# test/integration-tests/TEST-63-PATH/TEST-63-PATH.units/test63.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
ConditionPathExists=/tmp/nonexistent

[Service]
ExecStartPre=sh -c 'test "$TRIGGER_PATH" = /tmp/test63'
ExecStartPre=sh -c 'test "$TRIGGER_UNIT" = test63.path'
ExecStart=true
