# test/TEST-04-JOURNAL.units/logs-filtering.service
[Unit]
Description=Log filtering unit

[Service]
Type=oneshot
ExecStart=sh -c 'echo "Logging from the service, and ~more~ foo bar" && sleep 2'
SyslogLevel=notice
