# test/TEST-16-EXTEND-TIMEOUT.units/success-runtime.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Testsuite: Success Runtime (EXTEND_TIMEOUT_USEC > WATCHDOG_USEC however < RuntimeMaxSec)

[Service]

# EXTEND_TIMEOUT_USEC=4 second once during runtime, but sleep for 6 seconds.
# Runtime is 6 seconds and < RuntimeMaxSec so still successful.
Type=notify
TimeoutStartSec=4
TimeoutStopSec=4
RuntimeMaxSec=8
Environment=SERVICE=success_runtime extend_timeout_interval=4 sleep_interval=6 start_intervals=0 run_intervals=1 stop_intervals=0
ExecStart=/usr/lib/systemd/tests/testdata/TEST-16-EXTEND-TIMEOUT.units/extend-timeout.sh
