# mkosi.images/system/mkosi.extra/usr/lib/systemd/system/mnt.mount
[Mount]
What=/root/btrfs.raw
Where=/mnt
Type=btrfs
