# test/integration-tests/TEST-04-JOURNAL/TEST-04-JOURNAL.units/logs-filtering-syslog.service
[Unit]
Description=Log filtering unit

[Service]
Type=oneshot
ExecStart=/usr/lib/systemd/tests/testdata/TEST-04-JOURNAL.units/logs-filtering-syslog.py
SyslogLevel=notice
LogLevelMax=info
