# test/integration-tests/TEST-04-JOURNAL/TEST-04-JOURNAL.units/silent-success.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Silent successful service

[Service]
Type=oneshot
ExecStart=true
LogLevelMax=notice
