# mkosi/mkosi.images/minimal-0/mkosi.extra/usr/lib/systemd/user/minimal-app0.service
# SPDX-License-Identifier: LGPL-2.1-or-later

[Service]
ExecStartPre=cat /usr/lib/os-release
ExecStart=sleep 120
