# test/a-conj.service
[Unit]
Description=A conjugate
Requires=a.service
After=a.service
Before=a.service

[Service]
ExecStart=/bin/true
