# test/a.service
[Unit]
Description=A
Requires=b.service
Before=b.service

[Service]
ExecStart=/bin/true
