# test/dml-discard-empty.service
[Unit]
Description=DML discard empty service

[Service]
Slice=dml-discard.slice
Type=oneshot
ExecStart=/bin/true
