# test/dml-override.slice
[Unit]
Description=DML override slice

[Slice]
DefaultMemoryLow=10
