# test/dml-passthrough.slice
[Unit]
Description=DML passthrough slice

[Slice]
MemoryLow=100
