# test/dml-discard-set-ml.service
[Unit]
Description=DML discard set ml service

[Service]
Slice=dml-discard.slice
Type=oneshot
ExecStart=/bin/true
MemoryLow=15
