# test/dml.slice
[Unit]
Description=DML slice

[Slice]
DefaultMemoryLow=50
