# test/end.service.in
[Unit]
Description=End the test
After=testsuite.service

[Service]
ExecStart=@SYSTEMCTL@ poweroff --no-block
