# test/exec-group.service
[Unit]
Description=Test for Group

[Service]
ExecStart=/bin/sh -c 'exit $(test $(id -n -g) = nobody)'
Group=nobody
