# test/exec-user.service
[Unit]
Description=Test for User

[Service]
ExecStart=/bin/sh -c 'exit $(test "$USER" = nobody)'
User=nobody
