# test/exec-systemcallerrornumber.service
[Unit]
Description=Test for SystemCallErrorNumber

[Service]
ExecStart=/usr/bin/uname -a
SystemCallFilter=~uname
SystemCallErrorNumber=EACCES
