# test/exec-capabilityboundingset-invert.service
[Unit]
Description=Test for CapabilityBoundingSet

[Service]
ExecStart=/bin/sh -c 'c=$(capsh --print | grep "Bounding set " | cut -f 2 -d "="); echo $c; exit $(test -n $c)'
CapabilityBoundingSet=~
