# test/exec-umask-default.service
[Unit]
Description=Test for UMask default

[Service]
ExecStart=/bin/sh -c 'touch /tmp/test-exec-umask; s=$(stat -c %a /tmp/test-exec-umask); echo $s; exit $(test $s = "644")'
PrivateTmp=yes
