# test/exec-ignoresigpipe-no.service
[Unit]
Description=Test for IgnoreSIGPIPE=no

[Service]
ExecStart=/bin/sh -c 'kill -PIPE 0'
Type=oneshot
IgnoreSIGPIPE=no
