# test/hello.service
[Unit]
Description=Hello World

[Service]
ExecStart=/bin/echo "Hello World"
