# test/hello-after-sleep.target
[Unit]
Description=Sleep for a minute, then say hello.
Wants=sleep.service hello.service
After=sleep.service
Before=hello.service
