# test/nomem.slice
[Unit]
Description=Nomem Parent Slice

[Slice]
DisableControllers=memory
