# test/parent.slice
[Unit]
Description=Parent Slice

[Slice]
BlockIOWeight=200
