# test/parent-deep.slice
[Unit]
Description=Deeper Parent Slice

[Slice]
MemoryLimit=3G
