# test/parent.slice
[Unit]
Description=Parent Slice

[Slice]
IOWeight=200
