# test/path-mycustomunit.service
[Unit]
Description=Service Test Path Unit=

[Service]
ExecStart=/bin/true
Type=oneshot
