# test/path-unit.path
[Unit]
Description=Test Path Unit=

[Path]
PathExists=/tmp/test-path_unit
Unit=path-mycustomunit.service

[Install]
WantedBy=multi-user.target
