# test/path-changed.path
[Unit]
Description=Test PathChanged

[Path]
PathChanged=/tmp/test-path_changed

[Install]
WantedBy=multi-user.target
