# test/sched_idle_bad.service
[Unit]
Description=Bad sched priority for Idle

[Service]
ExecStart=/bin/true
CPUSchedulingPriority=1
