# test/test.service.in
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=%N
Wants=basic.target network.target systemd-user-sessions.service
After=basic.target network.target systemd-user-sessions.service
@unit@

[Service]
ExecStartPre=rm -f /failed /testok
ExecStart=@command@
Type=oneshot
MemoryAccounting=@memory-accounting@
Environment=SYSTEMD_PAGER= @env@
UnsetEnvironment=JOURNAL_STREAM
@service@
