# test/unstoppable.service
[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/echo 'I'm unstoppable!'
ExecStop=/bin/systemctl start --no-block unstoppable.service
