# test/testsuite-04.units/forever-print-hola.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=ForeverPrintHola service

[Service]
Type=exec
ExecStart=sh -x -c 'while :; do printf "Hola\n" || touch /i-lose-my-logs; sleep 1; done'
