# test/testsuite-04.units/logs-filtering.service
[Unit]
Description=Log filtering unit

[Service]
Type=oneshot
ExecStart=echo "Logging from the service, and ~more~ foo bar"
SyslogLevel=notice
