# test/testsuite-04.units/silent-success.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Silent successful service

[Service]
LogLevelMax=notice
ExecStart=/bin/true
