# test/testsuite-04.units/verbose-success.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Verbose successful service

[Service]
ExecStart=sh -c 'for fun in x y z ;do echo success ;done'
