# init/udev-control.socket
[Unit]
Description=udev Control Socket
DefaultDependencies=no
ConditionCapability=CAP_MKNOD

[Socket]
Service=udev.service
ListenSequentialPacket=/run/udev/control
