# init/udev.socket
[Unit]
Description=udev Kernel Device Manager Socket

[Socket]
ListenDatagram=@/org/kernel/udev/udevd
