# init/udev.socket
[Unit]
Description=udev Kernel Device Manager Socket
DefaultDependencies=no

[Socket]
ListenSequentialPacket=@/org/kernel/udev/udevd
