# init/udev.service.in
[Unit]
Description=udev Kernel Device Manager
Wants=udev.socket
After=udev.socket
Before=basic.target
DefaultDependencies=no

[Service]
Type=notify
OOMScoreAdjust=-1000
Restart=on-failure
ExecStart=@sbindir@/udevd
