# init/udev.service.in
[Unit]
Description=udev Kernel Device Manager
Requires=udev.socket
After=udev.socket
Before=basic.target
DefaultDependencies=no

[Service]
Type=notify
ExecStart=@sbindir@/udevd
