# init/udev-retry.service.in
[Unit]
Description=udev - re-trigger failed events after filesystems are available
DefaultDependencies=no
After=local-fs.target udev.service

[Service]
Type=oneshot
ExecStart=@sbindir@/udevadm trigger --type=failed --action=add

[Install]
WantedBy=sysinit.target
Also=udev.service udev-settle.service
