# init/udev-settle.service.in
[Unit]
Description=wait until full initialization
After=udev.service

[Service]
Type=finish
ValidNoProcess=yes
ExecStart=@sbindir@/udevadm settle

[Install]
WantedBy=sysinit.target
Also=udev.service
Also=udev-retry.service
