# init/udev-retry.service.in
[Unit]
Description=re-trigger failed events after filesystems are available
DefaultDependencies=no
After=local-fs.target udev.service

[Service]
Type=finish
ValidNoProcess=yes
ExecStart=@sbindir@/udevadm trigger --type=failed

[Install]
WantedBy=sysinit.target
Also=udev.service udev-settle.service
