# init/udev-settle.service.in
[Unit]
Description=wait until full initialization
DefaultDependencies=no
After=udev.service

[Service]
Type=oneshot
ValidNoProcess=yes
ExecStart=@sbindir@/udevadm settle

[Install]
WantedBy=sysinit.target
Also=udev.service udev-retry.service
