# init/udev-control.socket
[Unit]
Description=udev Control Socket
DefaultDependencies=no

[Socket]
Service=udev.service
ListenSequentialPacket=@/org/kernel/udev/udevd
