# init/udev.service.in
[Unit]
Description=/dev and kernel device manager
DefaultDependencies=no
Before=sysinit.target

[Service]
Type=notify
ExecStart=@sbindir@/udevd
ExecStartPost=@sbindir@/udevadm trigger --type=subsystems --action=add
ExecStartPost=@sbindir@/udevadm trigger --type=devices --action=add

[Install]
WantedBy=sysinit.target
Also=udev-settle.service
Also=udev-retry.service
