# init/udev-control.socket
[Unit]
Description=udev Control Socket
DefaultDependencies=no
ConditionVirtualization=!container

[Socket]
Service=udev.service
ListenSequentialPacket=/run/udev/control
