# test/fuzz/fuzz-unit-file/directives.path
path
[Path]
DirectoryMode=
DirectoryNotEmpty=
MakeDirectory=
PathChanged=
PathExists=
PathExistsGlob=
PathModified=
TriggerLimitBurst=
TriggerLimitIntervalSec=
Unit=
