# test/fuzz/fuzz-unit-file/syslog.socket
socket
#  SPDX-License-Identifier: LGPL-2.1-or-later
#
#  This file is part of systemd.
#
#  systemd is free software; you can redistribute it and/or modify it
#  under the terms of the GNU Lesser General Public License as published by
#  the Free Software Foundation; either version 2.1 of the License, or
#  (at your option) any later version.

[Unit]
Description=Syslog Socket
Documentation=man:systemd.special(7)
Documentation=https://systemd.io/SYSLOG
DefaultDependencies=no
Before=sockets.target

# Don't allow logging until the very end
Conflicts=shutdown.target
Before=shutdown.target

# Don't try to activate syslog.service if sysinit.target has failed.
Conflicts=emergency.service
Before=emergency.service

[Socket]
ListenDatagram=/run/systemd/journal/syslog
SocketMode=0666
PassCredentials=yes
PassSecurity=yes
ReceiveBuffer=8M

# The default syslog implementation should make syslog.service a
# symlink to itself, so that this socket activates the right actual
# syslog service.
#
# Examples:
#
# /etc/systemd/system/syslog.service -> /lib/systemd/system/rsyslog.service
# /etc/systemd/system/syslog.service -> /lib/systemd/system/syslog-ng.service
#
# Best way to achieve that is by adding this to your unit file
# (i.e. to rsyslog.service or syslog-ng.service):
#
# [Install]
# Alias=syslog.service
#
# See https://systemd.io/SYSLOG for details.

[Socket]
ListenStream=1.2.3.4:1234
ListenDatagram=1.2.3.4:1234
ListenSequentialPacket=1.2.3.4:1234
ListenFIFO=
ListenSpecial=
#ListenNetlink=
ListenMessageQueue=
ListenUSBFunction=
SocketProtocol=udplite
SocketProtocol=sctp
SocketProtocol=mptcp
SocketProtocol=
BindIPv6Only=false
Backlog=33
BindToDevice=eth0
SocketUser=daemon
SocketGroup=nobody
SocketMode=0111
DirectoryMode=0555
Accept=true
Accept=false
Writable=true
MaxConnections=11
MaxConnectionsPerSource=12
KeepAlive=yes
KeepAliveTimeSec=12345
KeepAliveIntervalSec=12345
KeepAliveProbes=12345
NoDelay=true
Priority=0
DeferAcceptSec=1
ReceiveBuffer=1G
SendBuffer=1G
IPTOS=low-delay
IPTOS=throughput
IPTOS=reliability
IPTOS=low-cost
IPTOS=
IPTTL=7
Mark=123
ReusePort=true
SmackLabel=smack-label
SmackLabelIPIn=smack-label
SmackLabelIPOut=no idea what to put here
SELinuxContextFromNet=true
PipeSize=11111
MessageQueueMaxMessages=200
MessageQueueMessageSize=200
FreeBind=false
Transparent=true
Broadcast=true
PassCredentials=true
PassSecurity=true
TCPCongestion=westwood
TCPCongestion=veno
TCPCongestion=cubic
TCPCongestion=lp
ExecStartPre=/bin/true "arg ' ' "
ExecStartPost=-!!/bin/false
ExecStopPre=/bin/true
ExecStopPost=-!!/bin/false
TimeoutSec=2343
Symlinks=a b c d e
Symlinks=
Symlinks=/a /b /c /d /e
FileDescriptorName=name
TriggerLimitIntervalSec=2343
TriggerLimitBurst=234
