# mkosi.images/minimal_1/mkosi.extra/usr/lib/systemd/system/minimal-app0.service
[Service]
ExecStartPre=cat /usr/lib/os-release
ExecStart=sleep 120
