# test/TEST-15-DROPIN/testsuite.service
[Unit]
Description=Testsuite service

[Service]
ExecStart=/test-dropin.sh
Type=oneshot
