# test/TEST-15-DROPIN/testsuite.service
[Unit]
Description=Testsuite service
After=multi-user.target

[Service]
ExecStart=/test-dropin.sh
Type=oneshot
