# test/test-fstab-generator/test-18-options.fstab.expected/mnt-rwonly.mount
# Automatically generated by systemd-fstab-generator

[Unit]
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
SourcePath=/etc/fstab
Before=local-fs.target
After=blockdev@dev-sdx11.target

[Mount]
What=/dev/sdx11
Where=/mnt/rwonly
Options=x-systemd.rw-only
ReadWriteOnly=yes
