# test/test-fstab-generator/test-18-options.fstab.expected/mnt-requiredby.mount
# Automatically generated by systemd-fstab-generator

[Unit]
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
SourcePath=/etc/fstab
DefaultDependencies=no
Conflicts=umount.target
Before=umount.target
After=blockdev@dev-sdx8.target

[Mount]
What=/dev/sdx8
Where=/mnt/requiredby
Options=x-systemd.required-by=foo.service
