# test/test-fstab-generator/test-18-options.fstab.expected/mnt-before.mount
# Automatically generated by systemd-fstab-generator

[Unit]
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
SourcePath=/etc/fstab
Before=foo.service
Before=local-fs.target
After=blockdev@dev-sdx4.target

[Mount]
What=/dev/sdx4
Where=/mnt/before
Options=x-systemd.before=foo.service
