# test/test-fstab-generator/test-18-options.fstab.expected/mnt-automount2.mount
# Automatically generated by systemd-fstab-generator

[Unit]
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
SourcePath=/etc/fstab
After=blockdev@dev-sdx10.target

[Mount]
What=/dev/sdx10
Where=/mnt/automount2
Options=x-systemd.automount,nofail
