# test/test-fstab-generator/test-18-options.fstab.expected/mnt-requires.mount
# Automatically generated by systemd-fstab-generator

[Unit]
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
SourcePath=/etc/fstab
After=foo.service
Requires=foo.service
Before=local-fs.target
After=blockdev@dev-sdx5.target

[Mount]
What=/dev/sdx5
Where=/mnt/requires
Options=x-systemd.requires=foo.service
