# test/test-fstab-generator/test-18-options.fstab.expected/mnt-onfailurex2dmulti.mount
# Automatically generated by systemd-fstab-generator

[Unit]
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
SourcePath=/etc/fstab
OnFailure=foo.service
OnFailure=bar.service
Before=local-fs.target
After=blockdev@dev-sdx19.target

[Mount]
What=/dev/sdx19
Where=/mnt/onfailure-multi
Options=x-systemd.onfailure=foo.service,x-systemd.onfailure=bar.service
