# test/test-fstab-generator/test-18-options.fstab.expected/mnt-nofail.mount
# Automatically generated by systemd-fstab-generator

[Unit]
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
SourcePath=/etc/fstab
After=blockdev@dev-sdx16.target

[Mount]
What=/dev/sdx16
Where=/mnt/nofail
Options=nofail
