# test/test-fstab-generator/test-18-options.fstab.expected/systemd-makefs@dev-sdx12.service
# Automatically generated by systemd-fstab-generator

[Unit]
Description=Make File System on %f
Documentation=man:systemd-makefs@.service(8)

DefaultDependencies=no
BindsTo=%i.device
After=%i.device
Before=systemd-fsck@%i.service mnt-mkfs.mount
Conflicts=shutdown.target
Before=shutdown.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/lib/systemd/systemd-makefs ext4 /dev/sdx12
TimeoutSec=infinity
