# test/TEST-03-JOBS/TEST-03-JOBS.units/hello.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Hello World

[Service]
ExecStart=/bin/echo "Hello World"
