# test/TEST-03-JOBS/TEST-03-JOBS.units/shortcut-restart.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Service]
Type=oneshot
ExecStartPre=sleep 10
ExecStart=false
Restart=on-failure
RestartSec=1d
