# test/TEST-03-JOBS/TEST-03-JOBS.units/sleep-infinity-restart-normal.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
OnFailure=restart-on-failure.service

[Service]
ExecStart=/bin/sleep infinity
Restart=on-failure
RestartMode=normal
