# units/session/exit.service.in
[Unit]
Description=Exit
DefaultDependencies=no
Requires=shutdown.target
After=shutdown.target

[Service]
Type=oneshot
ExecStart=@SYSTEMCTL@ --session daemon-exit
