# postfix.socket
[Meta]
Description=Postfix SMTP Socket

[Socket]
Listen=25
Type=stream
