# syslog.socket
[Meta]
Description=Syslog Socket

[Socket]
Listen=/dev/log
Type=dgram
