# test/testsuite-10.units/test10.service
[Unit]
Requires=test10.socket
ConditionPathExistsGlob=/tmp/nonexistent

[Service]
ExecStart=true
