# units/fedora/sysinit.service
#  This file is part of systemd.
#
#  systemd is free software; you can redistribute it and/or modify it
#  under the terms of the GNU General Public License as published by
#  the Free Software Foundation; either version 2 of the License, or
#  (at your option) any later version.

[Unit]
Description=System Initialization
DefaultDependencies=no
Before=local-fs.target remote-fs.target swap.target sockets.target
Conflicts=shutdown.target
Before=shutdown.target

[Service]
ExecStart=/etc/rc.d/rc.sysinit
Type=forking
TimeoutSec=0
StandardInput=tty
ValidNoProcess=yes
