# mkosi/mkosi.initrd.conf/mkosi.extra/usr/lib/systemd/system/systemd-nspawn@.service.d/fdstore.conf
[Service]
FileDescriptorStoreMax=16
FileDescriptorStorePreserve=yes
