# test/integration-tests/TEST-54-CREDS/systemd.unit-dropin.my-service.service~30-named.service
# SPDX-License-Identifier: LGPL-2.1-or-later

[Service]
ExecStart=touch /tmp/unit-named-dropin
