# mkosi.images/testsuite-69-shutdown/mkosi.extra/usr/lib/systemd/system/systemd-firstboot.service.d/no-root-password.conf
[Service]
ExecStart=systemd-firstboot --prompt-locale --prompt-timezone --delete-root-password
