# test/TEST-07-PID1/TEST-07-PID1.units/issue27953.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=This unit should not remain active once the shell process exits

[Service]
Type=oneshot
ExecStart=sh -c 'sleep infinity & exit 0'
