# test/TEST-23-UNIT-FILE/TEST-23-UNIT-FILE.units/TEST-23-UNIT-FILE-joins-namespace-of-3.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
JoinsNamespaceOf=TEST-23-UNIT-FILE-joins-namespace-of-1.service

[Service]
Type=oneshot
MountAPIVFS=yes
PrivateTmp=yes
ExecStart=test -e /tmp/shared-private-file
ExecStart=test -e /tmp/hoge
