# test/test-cgroup-unit-default/dml-discard.slice
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=DML discard slice

[Slice]
DefaultMemoryLow=
