# test/test-cgroup-unit-default/dml.slice
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=DML slice

[Slice]
DefaultMemoryLow=50
