# mkosi.sanitizers/mkosi.extra/usr/lib/systemd/system/systemd-journald.service.d/10-stdout-tty.conf
# SPDX-License-Identifier: LGPL-2.1-or-later

[Service]
# Sanitizers log to stderr by default. However, journald's stderr is connected
# to /dev/null, so we lose all the sanitizer logs. To rectify that, let's
# connect journald's stdout to kmsg so that the sanitizer failures end up in the
# journal.
StandardOutput=kmsg
