# units/user/systemd-ptybrokerd.service.in
#  SPDX-License-Identifier: LGPL-2.1-or-later
#
#  This file is part of systemd.
#
#  systemd is free software; you can redistribute it and/or modify it
#  under the terms of the GNU Lesser General Public License as published by
#  the Free Software Foundation; either version 2.1 of the License, or
#  (at your option) any later version.

[Unit]
Description=PTY Broker
Documentation=man:systemd-ptybrokerd.service(8)
Requires=systemd-ptybrokerd.socket
After=systemd-ptybrokerd.socket

[Service]
ExecStart={{LIBEXECDIR}}/systemd-ptybrokerd --user
Type=notify
Restart=on-failure
RestartSec=0
{{SERVICE_WATCHDOG}}

[Install]
Also=systemd-ptybrokerd.socket
