# units/user/systemd-journalctl@.service
[Unit]
Description=Journal Log Access Service
Documentation=man:journalctl(1)

[Service]
ExecStart=journalctl
