# test/integration-test/TEST-07-PID1/TEST-07-PID1.units/issue14566-repro.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Issue 14566 Repro

[Service]
ExecStart=/usr/lib/systemd/tests/testdata/TEST-07-PID1.units/%N.sh
ExecStopPost=/bin/true
KillMode=mixed
