# test1/exec-demo.service
[Meta]
Description=Simple Execution Demo

[Service]
ExecStart=/bin/cat /etc/hosts
Type=simple
Output=syslog
