# test1/fifo-demo.service
[Meta]
Description=FIFO Activation Demo Service

[Service]
ExecStartPre=/bin/echo "Wow, es geht los!"
#ExecStart=/bin/bash -c 'read a <&3 ; echo "GOT: $a"'
ExecStart=/bin/bash -c 'cat <&3'
Type=simple
