# test1/postfix.socket
[Meta]
Description=Postfix SMTP Socket

[Socket]
ListenStream=53333
ListenFIFO=/tmp/systemd-fifo
