# test1/postfix.socket
[Meta]
Description=Postfix SMTP Socket

[Socket]
ListenStream=25
ListenFIFO=/dev/test
