# test1/syslog.socket
[Meta]
Description=Syslog Socket

[Socket]
ListenDatagram=/dev/log
ListenStream=eth0:4711
