# test1/syslog.service
[Meta]
Description=System Logging Daemon

[Service]
ExecStart=/usr/bin/rsyslogd --foobar 'waldo'
