# test1/syslog.socket
[Meta]
Description=Syslog Socket

[Socket]
ListenDatagram=/tmp/systemd-socket
ListenStream=eth0:3456
