# test1/syslog.service
[Meta]
Description=System Logging Daemon

[Service]
Exec=/usr/bin/rsyslogd
