# test1/systemd-logger.socket
[Meta]
Description=systemd Logging Socket

[Socket]
ListenStream=@/systemd/logger
