# test1/systemd-logger.service
[Meta]
Description=systemd Logging Daemon

[Service]
ExecStart=/home/lennart/projects/systemd/systemd-logger
Type=simple
