# test/test-fstab-generator/test-14-uuid.expected/sysroot.mount
# Automatically generated by systemd-fstab-generator

[Unit]
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
SourcePath=/proc/cmdline
Before=initrd-root-fs.target
Requires=systemd-fsck-root.service
After=systemd-fsck-root.service
After=blockdev@dev-disk-by\x2duuid-3f5ad593\x2d4546\x2d4a94\x2da374\x2dbcfb68aa11f7.target

[Mount]
What=/dev/disk/by-uuid/3f5ad593-4546-4a94-a374-bcfb68aa11f7
Where=/sysroot
Options=ro
