# mkosi.extra/usr/lib/systemd/system/iscsi-init.service.d/asan.conf
# SPDX-License-Identifier: LGPL-2.1-or-later

# The iscsi-init.service calls `sh` which might, in certain circumstances, pull in instrumented systemd NSS
# modules causing `sh` to fail. Avoid the issue by setting LD_PRELOAD to load the sanitizer libraries if
# needed.
[Service]
EnvironmentFile=-/usr/lib/systemd/systemd-asan-env
