# test/testsuite-83.units/mock-polkit.service
[Service]
Type=dbus
BusName=org.freedesktop.PolicyKit1
ExecStart=/usr/lib/systemd/tests/unit-tests/manual/test-bus-polkit --mock-polkit
