# test/testsuite-83.units/test-bus-polkit.service
[Service]
Type=dbus
BusName=io.systemd.test.TestBusPolkit
ExecStart=/usr/lib/systemd/tests/unit-tests/manual/test-bus-polkit
