# test/test-execute/exec-personality-x86-64.service
[Unit]
Description=Test for Personality=x86-64

[Service]
ExecStart=/bin/sh -c 'echo $(uname -m); exit $(test $(uname -m) = "x86_64")'
Type=oneshot
Personality=x86-64
