# test/test-execute/exec-ioschedulingclass-idle.service
[Unit]
Description=Test for IOSchedulingClass=idle

[Service]
ExecStart=/bin/bash -c 'c=$(ionice); echo $c; [[ "$c" == idle* ]]'
Type=oneshot
IOSchedulingClass=idle
