# test/test-execute/exec-restrictnamespaces-merge-and.service
[Unit]
Description=Test merging RestrictNamespaces= with AND

[Service]
RestrictNamespaces=mnt pid
RestrictNamespaces=~mnt usr
ExecStart=unshare -p -f
ExecStart=sh -c '! unshare -m -u -i -n'
Type=oneshot
