# test/test-execute/exec-supplementarygroups.service
[Unit]
Description=Test for Supplementary Group

[Service]
ExecStart=/bin/sh -x -c 'test " $$(id -G) " = *" 0 1 "*'
Type=oneshot
SupplementaryGroups=1
