# test/test-execute/exec-privatetmp-disconnected-nodefaultdeps-nor-sandboxing.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Test for PrivateTmp=disconnected with DefaultDependencies=no + no sandboxing
DefaultDependencies=no

[Service]
ExecStart=+touch /tmp/test-exec_privatetmp_disconnected
ExecStart=+sh -x -c 'test x$$TMPDIR = x'
Type=oneshot
PrivateTmp=disconnected
