# test/test-execute/exec-binarydirectory.service
[Service]
ExecStart=ls_temp
Type=oneshot
BinaryDirectory=/tmp/test-exec_binarydirectory
