# test/test-execute/exec-systemcallfilter-failing2.service
[Unit]
Description=Test for SystemCallFilter

[Service]
ExecStart=/bin/echo "This should not be seen"
Type=oneshot
SystemCallFilter=~write open execve exit_group close mmap munmap fstat DONOTEXIST
