# test/test-execute/exec-privatedevices-yes.service
[Unit]
Description=Test for PrivateDev=yes

[Service]
ExecStart=/bin/sh -c 'test ! -c /dev/kmsg'
Type=oneshot
PrivateDevices=yes
