# test/test-execute/exec-ignoresigpipe-no.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Test for IgnoreSIGPIPE=no

[Service]
ExecStart=/bin/sh -x -c 'kill -PIPE 0'
Type=oneshot
IgnoreSIGPIPE=no
