# test/test-execute/exec-execsearchpaths.service
[Service]
ExecStart=ls_temp
Type=oneshot
ExecSearchPaths=/does/not/exist /tmp/test-exec_execsearchpaths
