# test/test-execute/exec-user.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Test for User (daemon)

[Service]
ExecStart=sh -x -c 'test "$$USER" = "daemon"'
Type=oneshot
User=daemon
