# test/test-execute/exec-restrictnamespaces-yes.service
[Unit]
Description=Test RestrictNamespaces=yes

[Service]
RestrictNamespaces=yes
ExecStart=unshare -m
Type=oneshot
