# test/test-execute/exec-privatenetwork-yes.service
[Unit]
Description=Test for PrivateNetwork

[Service]
ExecStart=/bin/sh -x -c 'i=$$(ip link | grep ": " | grep -v ": lo:"); test -z "$$i"'
Type=oneshot
PrivateNetwork=yes
