# test/test-execute/exec-runtimedirectory-mode.service
[Unit]
Description=Test for RuntimeDirectoryMode

[Service]
ExecStart=/bin/sh -c 's=$(stat -c %a /tmp/test-exec_runtimedirectory-mode); echo $s; exit $(test $s = "750")'
Type=oneshot
RuntimeDirectory=test-exec_runtimedirectory-mode
RuntimeDirectoryMode=0750
