# test/test-execute/exec-ioschedulingclass-none.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Test for IOSchedulingClass=none

[Service]
# Old kernels might report "none" here, new kernels "best-effort".
ExecStart=/bin/sh -x -c 'c=$$(LC_ALL=C ionice); test "$${c%%:*}" = "none" -o "$${c%%:*}" = "best-effort"'
Type=oneshot
IOSchedulingClass=none
