# test/test-execute/exec-credentials-dir-specifier.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Test for specifiers

[Service]
Type=oneshot
Environment=TOP_SECRET=%d/very_top_secret
# Test if the specifier is resolved correctly both before and after LoadCredential=
ExecStart=test %d/very_top_secret = "${CREDENTIALS_DIRECTORY}/very_top_secret"
LoadCredential=very_top_secret
ExecStart=test %d/very_top_secret = "${CREDENTIALS_DIRECTORY}/very_top_secret"
ExecStart=sh -c 'test %d/very_top_secret = "$TOP_SECRET"'
