# test/test-execute/exec-restrictnamespaces-merge-all.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Test merging RestrictNamespaces= with all flags

[Service]
RestrictNamespaces=mnt pid cgroup net uts ipc user
RestrictNamespaces=net
ExecStart=unshare -m -u -i -n -p -f
Type=oneshot
