# test/test-execute/exec-systemcallerrornumber.service
[Unit]
Description=Test for SystemCallErrorNumber

[Service]
ExecStart=/bin/sh -c 'uname -a'
SystemCallFilter=~uname
SystemCallErrorNumber=EACCES
