# test/test-execute/exec-systemcallerrornumber.service
[Unit]
Description=Test for SystemCallErrorNumber

[Service]
ExecStart=/bin/sh -x -c 'uname -a'
Type=oneshot
SystemCallFilter=~uname
SystemCallErrorNumber=EACCES
