# test/test-execute/exec-protectkernelmodules-no-capabilities.service
[Unit]
Description=Test CAP_SYS_MODULE ProtectKernelModules=no

[Service]
ProtectKernelModules=no
# sed: remove dropped capabilities (cap_xxx-[epi]) from the output
ExecStart=/bin/sh -x -c 'capsh --print | sed -r "s/[^ ]+?\-[epi]+//g" | grep cap_sys_module'
Type=oneshot
