# test/test-execute/exec-privatedevices-yes-capability-mknod.service
[Unit]
Description=Test CAP_MKNOD capability for PrivateDevices=yes

[Service]
PrivateDevices=yes
ExecStart=/bin/sh -x -c '! capsh --print | grep cap_mknod'
Type=oneshot
