# test/test-execute/exec-ioschedulingclass-realtime.service
[Unit]
Description=Test for IOSchedulingClass=realtime

[Service]
ExecStart=/bin/bash -c 'c=$(ionice); echo $c; [[ "$c" == realtime* ]]'
Type=oneshot
IOSchedulingClass=realtime
