# test/test-execute/exec-protectkernellogs-yes-capabilities.service
[Unit]
Description=Test CAP_SYSLOG for ProtectKernelLogs=yes

[Service]
ProtectKernelLogs=yes
ExecStart=/bin/sh -x -c '! capsh --print | grep cap_syslog'
Type=oneshot
