# test/test-execute/exec-passenvironment-empty.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Test for PassEnvironment and erasing the variable list

[Service]
ExecStart=/bin/sh -x -c 'test "$${VAR1-unset}" = "unset" && test "$${VAR2-unset}" = "unset" && test "$${VAR3-unset}" = "unset" && test "$${VAR4-unset}" = "unset" && test "$${VAR5-unset}" = "unset"'
Type=oneshot
PassEnvironment=VAR1 VAR2 VAR3 VAR4 VAR5
PassEnvironment=
