# test/test-execute/exec-protectkernelmodules-strict-capabilities.service
[Unit]
Description=Test CAP_SYS_MODULE for ProtectKernelModules=strict

[Service]
ProtectKernelModules=strict
ExecStart=/bin/sh -x -c '! capsh --print | grep cap_sys_module'
Type=oneshot
