# test/test-execute/exec-execsearchpath.service
[Service]
ExecStart=ls_temp
Type=oneshot
ExecSearchPath=/tmp/test-exec_execsearchpath
