# test/test-execute/exec-standardinput-data.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Test for StandardInputText= and StandardInputData=

[Service]
ExecStart=/bin/sh -x -c 'd=$$(mktemp -d -p /tmp); echo -e "this is a test\nand this is more\nsomething encoded!\nsomething   in multiple lines\nand some more\nand a more bas64 data\nsomething with strange\nembedded\tcharacters\nand something with a exec-stdin-data.service specifier" > $d/text ; cmp $d/text ; rm -rf $d'
Type=oneshot
StandardInput=data
StandardInputText=this is a test
StandardInputText=and this is more
StandardInputData=c29tZXRoaW5nIGVuY29kZWQhCg==
StandardInputText=something \
 in multiple lines
StandardInputText=\
and some more
StandardInputData=YW5kIGEgbW9y \
    ZSBiYXM2NCBk\
YXRhCg==
StandardInputText=something with strange\nembedded\tcharacters
StandardInputText=and something with a %n specifier
