# test/test-execute/exec-group.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Test for Group (daemon)

[Service]
ExecStart=/bin/sh -x -c 'test "$$(id -n -g)" = "daemon"'
Type=oneshot
Group=daemon
