# test/test-execute/exec-cpuaffinity1.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Test for CPUAffinity (simple)

[Service]
ExecStart=sh -c 'test $$(cat /proc/self/status | grep Cpus_allowed: | rev | cut -c 1) = 1'
CPUAffinity=0
