# test/test-execute/exec-personality-x86.service
[Unit]
Description=Test for Personality=x86

[Service]
ExecStart=/bin/sh -x -c 'c=$$(uname -m); test "$$c" = "i686"'
Type=oneshot
Personality=x86
