# test/test-execute/exec-privatetmp-no.service
[Unit]
Description=Test for PrivateTmp=no

[Service]
ExecStart=/bin/sh -x -c 'test -f /tmp/test-exec_privatetmp'
Type=oneshot
PrivateTmp=no
