# test/test-execute/exec-readonlypaths-mount-propagation.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Test to make sure that passing ReadOnlyPaths= disconnect mount propagation

[Service]
ReadOnlyPaths=-/i-dont-exist
ExecStart=sh -x -c 'd=$$(mktemp -d -p /tmp); trap "umount \'$$d\' && rmdir \'$$d\'" EXIT; mount -t tmpfs tmpfs "$$d"; grep "$$d" /proc/self/mountinfo && ! grep "$$d" /proc/$${PPID}/mountinfo && ! grep "$$d" /proc/1/mountinfo'
Type=oneshot
