# test/test-execute/exec-personality-x86-64.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Test for Personality=x86-64

[Service]
ExecStart=/bin/sh -x -c 'c=$$(uname -m); test "$$c" = "x86_64"'
Type=oneshot
Personality=x86-64
