# test/test-execute/exec-runtimedirectory-owner.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Test for RuntimeDirectory owner (must not be the default group of the user if Group is set)

[Service]
ExecStart=sh -x -c 'group=$$(stat -c %%G %t/test-exec_runtimedirectory-owner-daemon); test "$$group" = "daemon"'
Type=oneshot
Group=daemon
User=root
RuntimeDirectory=test-exec_runtimedirectory-owner-daemon
