# test/test-execute/exec-passenvironment.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Test for PassEnvironment

[Service]
ExecStart=sh -x -c 'test "$$VAR1" = "word1 word2" && test "$$VAR2" = word3 && test "$$VAR3" = "\\$$word 5 6" && test "$$VAR4" = "new\nline" && test "$$VAR5" = passwordwithbackslashes'
Type=oneshot
PassEnvironment=VAR1 VAR2 VAR3 VAR4 VAR5
