# test/test-execute/exec-execsearchpaths.service
[Service]
ExecStart=ls_temp
Type=oneshot
ExecSearchPaths=/tmp/test-exec_execsearchpaths
