# test/test-execute/exec-group-nogroup.service
[Unit]
Description=Test for Group

[Service]
ExecStart=/bin/sh -x -c 'test "$$(id -n -g)" = "nogroup"'
Type=oneshot
Group=nogroup
