# test/test-execute/exec-workingdirectory-trailing-dot.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Test for WorkingDirectory with trailing dot

[Service]
ExecStart=sh -x -c 'test "$$PWD" = "/tmp/test-exec_workingdirectory"'
Type=oneshot
WorkingDirectory=/tmp///./test-exec_workingdirectory/.
