# test/test-execute/exec-personality-ppc64.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Test for Personality=ppc64

[Service]
ExecStart=/bin/sh -c 'echo $(uname -m); exit $(test $(uname -m) = "ppc64")'
Type=oneshot
Personality=ppc64
