# test/test-execute/exec-user-nobody.service
[Unit]
Description=Test for User

[Service]
ExecStart=/bin/sh -x -c 'test "$$USER" = "nobody"'
Type=oneshot
User=nobody
