# test/test-execute/exec-specifier-interpolation.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=https://github.com/systemd/systemd/issues/2637

[Service]
Type=oneshot
ExecStart=bash -x -c "[[ %%U == ?U ]]"
