# test/test-execute/exec-ignoresigpipe-yes.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Test for IgnoreSIGPIPE=yes

[Service]
ExecStart=sh -x -c 'kill -PIPE 0'
Type=oneshot
IgnoreSIGPIPE=yes
