# test/test-execute/exec-dynamicuser-fixeduser-one-supplementarygroup.service
[Unit]
Description=Test DynamicUser with User= and SupplementaryGroups=

[Service]
ExecStart=/bin/sh -x -c 'HAVE=0; for g in $$(id -G); do test "$$g" = "1" && HAVE=1; done; test "$$HAVE" -eq 1'
ExecStart=/bin/sh -x -c 'test "$$(id -g)" = "1" && test "$$(id -u)" = "1"'
Type=oneshot
User=1
DynamicUser=yes
SupplementaryGroups=1
