# test/test-execute/exec-supplementarygroups-single-group-user.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Test for Supplementary Group with only one group and uid 1

[Service]
ExecStart=/bin/sh -x -c 'for g in $$(id -G); do test "$$g" = "1" && exit 0; done; exit 1'
ExecStart=/bin/sh -x -c 'test "$$(id -g)" = "1" && test "$$(id -u)" = "1"'
Type=oneshot
User=1
Group=1
SupplementaryGroups=1
