# test/test-execute/exec-privatetmp-yes.service
[Unit]
Description=Test for PrivateTmp=yes

[Service]
ExecStart=/bin/sh -x -c 'test ! -f /tmp/test-exec_privatetmp'
Type=oneshot
PrivateTmp=yes
