# test/test-execute/exec-privatedevices-no.service
[Unit]
Description=Test for PrivateDev=no

[Service]
ExecStart=/bin/sh -x -c 'test -c /dev/mem'
Type=oneshot
PrivateDevices=no
