# test/test-execute/exec-restrictnamespaces-mnt-blacklist.service
[Unit]
Description=Test RestrictNamespaces=~mnt

[Service]
RestrictNamespaces=~mnt
ExecStart=unshare -m
Type=oneshot
