# test/test-execute/exec-readonlypaths-mount-propagation.service
[Unit]
Description=Test to make sure that passing ReadOnlyPaths= disconnect mount propagation

[Service]
ReadOnlyPaths=-/i-dont-exist
ExecStart=/bin/sh -x -c 'mkdir -p /TEST; mount -t tmpfs tmpfs /TEST; grep TEST /proc/self/mountinfo && ! grep TEST /proc/$${PPID}/mountinfo && ! grep TEST /proc/1/mountinfo'
Type=oneshot
