# test/test-execute/exec-personality-ppc64le.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Test for Personality=ppc64le

[Service]
ExecStart=sh -c 'echo $(uname -m); exit $(test $(uname -m) = "ppc64le")'
Type=oneshot
Personality=ppc64le
