# test/test-execute/exec-protectkernellogs-yes-capabilities.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Test CAP_SYSLOG for ProtectKernelLogs=yes

[Service]
ProtectKernelLogs=yes
# sed: remove dropped (cap_xxx-[epi]) and IAB capabilities from the output
ExecStart=sh -x -c '! capsh --print | sed -re "s/[^ ]+?\-[epi]+//g" -e '/IAB/d' | grep cap_syslog'
Type=oneshot
