# test/units/loopy.service.d/compat.conf
[Unit]
BindsTo=loopy2.service

[Install]
Also=loopy2.service
