# etc/udev/gentoo/90-network.rules
# do not edit this file, it will be overwritten on update

# /etc/udev/rules/90-network.rules:  triggering network init-scripts

# Activate our network if we can
SUBSYSTEM=="net", ACTION=="add",    RUN+="net.sh %k start"
SUBSYSTEM=="net", ACTION=="remove", RUN+="net.sh %k stop"

