# extras/udev-acl/70-acl.rules
# do not edit this file, it will be overwritten on update

ENV{MAJOR}=="", GOTO="acl_end"
ENV{ACL_MANAGE}=="0", GOTO="acl_end"
ACTION!="add|change", GOTO="acl_apply"

# PTP/MTP protocol devices, cameras, portable media players
SUBSYSTEM=="usb", ENV{ID_USB_INTERFACES}=="", ENV{DEVTYPE}=="usb_device", IMPORT{program}="usb_id --export %p"
SUBSYSTEM=="usb", ENV{ID_USB_INTERFACES}=="*:060101:*", ENV{ACL_MANAGE}="1"

# digicams with proprietary protocol
ENV{ID_GPHOTO2}=="*?", ENV{ACL_MANAGE}="1"

# SCSI scanners
KERNEL=="sg[0-9]*", ATTRS{type}=="6", ENV{ACL_MANAGE}="1"
KERNEL=="sg[0-9]*", ATTRS{type}=="3", ATTRS{vendor}=="HP|EPSON|Epson", ENV{ACL_MANAGE}="1"

# USB scanners
ENV{libsane_matched}=="yes", ENV{ACL_MANAGE}="1"

# HPLIP devices (necessary for ink level check and HP tool maintenance)
ENV{ID_HPLIP}=="1", ENV{ACL_MANAGE}="1"

# optical drives
SUBSYSTEM=="block", ENV{ID_CDROM}=="1", ENV{ACL_MANAGE}="1"

# sound devices
SUBSYSTEM=="sound", ENV{ACL_MANAGE}="1"
# sound jack-sense
SUBSYSTEM=="input", SUBSYSTEMS=="sound", ENV{ACL_MANAGE}="1"

# webcams, frame grabber, TV cards
SUBSYSTEM=="video4linux", ENV{ACL_MANAGE}="1"
SUBSYSTEM=="dvb", ENV{ACL_MANAGE}="1"

# fingerprint readers
SUBSYSTEM=="usb", ATTR{idVendor}=="0483", ATTR{idProduct}=="2016", ENV{ACL_MANAGE}="1"

# DRI video devices
SUBSYSTEM=="drm", KERNEL=="card*", ENV{ACL_MANAGE}="1"

# KVM
SUBSYSTEM=="misc", KERNEL=="kvm", ENV{ACL_MANAGE}="1"

# smart-card readers
ENV{ID_SMARTCARD_READER}=="*?", ENV{ACL_MANAGE}="1"

# joysticks
SUBSYSTEM=="input", ENV{ID_CLASS}=="joystick", ENV{ACL_MANAGE}="1"

# apply ACL for all locally logged in users
LABEL="acl_apply", ENV{ACL_MANAGE}=="?*", TEST=="/var/run/ConsoleKit/database", \
  RUN+="udev-acl --action=$env{ACTION} --device=$env{DEVNAME}"

LABEL="acl_end"
