# test/integration-tests/TEST-55-OOMD/TEST-55-OOMD.units/TEST-55-OOMD-testbloat.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Create a lot of memory pressure

[Service]
# A VERY small memory.high will cause the 'stress-ng' (trying to use a lot of memory)
# to throttle and be put under heavy pressure.
MemoryHigh=3M
Slice=TEST-55-OOMD-workload.slice
# Pin --vm-method=lfsr32: the only stress-ng vm method without TARGET_CLONES,
# so it can't dispatch to an AVX-512 clone and SIGILL on CPUs lacking it.
ExecStart=stress-ng --timeout 3m --vm 10 --vm-bytes 200M --vm-keep --vm-method=lfsr32
