# test/integration-tests/TEST-55-OOMD/TEST-55-OOMD.units/TEST-55-OOMD-testmunch.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Create some memory pressure

[Service]
MemoryHigh=12M
Slice=TEST-55-OOMD-workload.slice
# Pin --vm-method=lfsr32: the only stress-ng vm method without TARGET_CLONES,
# so it can't dispatch to an AVX-512 clone and SIGILL on CPUs lacking it.
ExecStart=stress-ng --timeout 3m --vm 10 --vm-bytes 200M --vm-keep --vm-method=lfsr32
