# mkosi.images/testsuite-54-creds-system/mkosi.extra/usr/lib/systemd/system/testsuite-35.service.d/99-stdout.conf
[Service]
StandardOutput=journal+console
StandardError=journal+console