# test/testsuite-63.units/test63.service
[Unit]
ConditionPathExists=!/tmp/nonexistent

[Service]
ExecStart=true
