# test/testsuite-63.units/test63-issue-24577-dep.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Service]
Type=oneshot
ExecStart=bash -c 'sleep infinity'
