# mkosi.images/system/mkosi.extra/usr/lib/systemd/system/user@.service.d/99-SYSTEMD_UNIT_PATH.conf
[Service]
PassEnvironment=SYSTEMD_UNIT_PATH
