# test/testsuite-30.units/systemd-timedated.service.d/watchdog.conf
# SPDX-License-Identifier: LGPL-2.1-or-later

[Service]
WatchdogSec=10min
