# test/TEST-04-JOURNAL/TEST-04-JOURNAL.units/verbose-success.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Verbose successful service

[Service]
Type=oneshot
# Sleep so that the cgroup is still there when journald processes the log message which is required for
# journald to add the expected fields to the log message.
ExecStart=sleep 2
ExecStart=echo success
ExecStart=sleep 2
