# test/test-fstab-generator/test-18-options.fstab.expected.sysroot/mnt-reqmounts.mount
# Automatically generated by systemd-fstab-generator

[Unit]
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
SourcePath=/etc/fstab
RequiresMountsFor=/hoge
Before=local-fs.target
After=blockdev@dev-sdx6.target

[Mount]
What=/dev/sdx6
Where=/mnt/reqmounts
Options=x-systemd.requires-mounts-for=/hoge
