# test/test-fstab-generator/test-18-options.fstab.expected.sysroot/mnt-wantedbyx2dautomount.automount
# Automatically generated by systemd-fstab-generator

[Unit]
SourcePath=/etc/fstab
Documentation=man:fstab(5) man:systemd-fstab-generator(8)

[Automount]
Where=/mnt/wantedby-automount
