# mkosi.images/testsuite-07-pid1-system/mkosi.extra/usr/lib/systemd/system/issue2730.mount
[Mount]
What=tmpfs
Where=/issue2730
Type=tmpfs

[Install]
WantedBy=local-fs.target
Alias=issue2730-alias.mount
