# test/TEST-47-HONORFIRSTSHUTDOWN/testsuite.service
[Unit]
Description=Testsuite service

[Service]
ExecStart=/testsuite.sh
Type=oneshot
