# test/integration-tests/TEST-52-HONORFIRSTSHUTDOWN/TEST-52-HONORFIRSTSHUTDOWN.units/test-honor-first-shutdown.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Honor First Shutdown feature
After=multi-user.target

[Service]
ExecStart=/usr/lib/systemd/tests/testdata/TEST-52-HONORFIRSTSHUTDOWN.units/%N.sh
ExecStop=sh -c 'kill -KILL $MAINPID'
FailureAction=reboot

[Install]
WantedBy=multi-user.target
