# units/user@0.service.d/10-defaults.conf
[Unit]
After=user-runtime-dir@0.service dbus.service
