# mkosi/mkosi.extra/usr/lib/systemd/system/user@.service.d/fdstore.conf
# For tests exercising the fd store we need the unit in the rootfs to have these
# settings, or the fdstore content will be dropped in the initrd -> rootfs transition
[Service]
FileDescriptorStoreMax=20
FileDescriptorStorePreserve=yes
