# test/test-socket/conf/test-protocol-udplite.socket
[Unit]
Description=Test Socket ListenStream SocketProtocol udplite

[Socket]
ListenDatagram=9998
SocketProtocol=udplite

[Install]
WantedBy=sockets.target
