# test/test-socket/conf/test-socket-datagram.socket
[Unit]
Description=Test Socket ListenDatagram

[Socket]
ListenDatagram=9999

[Install]
WantedBy=sockets.target
