# test/test-socket/conf/test-socket-stream-ipv6only.socket
[Unit]
Description=Test Socket ListenStream BindIPv6Only

[Socket]
ListenStream=9990
BindIPv6Only=both

[Install]
WantedBy=sockets.target
