# test/test-socket/conf/test-socket-netlink.socket
[Unit]
Description=Test Socket ListenNetlink

[Socket]
ListenNetlink=route 1361

[Install]
WantedBy=sockets.target
