# test/test-socket/conf/test.service
[Unit]
Description=Test socket Service

[Service]
Type=simple
ExecStart=/usr/bin/true
StandardInput=socket
StandardError=journal
TimeoutStopSec=5

[Install]
WantedBy=multi-user.target
