# test/test-socket/conf/test-protocol-sctp.socket
[Unit]
Description=Test Socket ListenStream SocketProtocol sctp

[Socket]
ListenStream=9998
SocketProtocol=sctp

[Install]
WantedBy=sockets.target
