# test/test-socket/conf/test-socket-vsock.socket
[Unit]
Description=Test Socket ListenStream=vsock:2:1234

[Socket]
ListenStream=vsock:2:1234

[Install]
WantedBy=sockets.target
