# test/test-socket/conf/test-socket-stream.socket
[Unit]
Description=Test Socket ListenStream

[Socket]
ListenStream=9999

[Install]
WantedBy=sockets.target
