# test/testsuite-81.units/cryptenroll.service
[Service]
Type=oneshot
ExecStart=/usr/lib/systemd/tests/testdata/testsuite-81.units/test.sh
