# test/integration-test/TEST-80-NOTIFYACCESS/TEST-80-NOTIFYACCESS.units/fdstore-pin.target
[Unit]
After=fdstore-pin.service fdstore-nopin.service
Wants=fdstore-pin.service fdstore-nopin.service
