# etc/udev/debian/run.rules
# debugging monitor
RUN+="socket:/org/kernel/udev/monitor"

# run a command on remove events
ACTION=="remove", ENV{REMOVE_CMD}!="", RUN+="$env{REMOVE_CMD}"

# ignore the events generated by virtual consoles
KERNEL=="ptmx",		OPTIONS+="last_rule"
KERNEL=="console",	OPTIONS+="last_rule"
KERNEL=="tty"	,	OPTIONS+="last_rule"
KERNEL=="tty[0-9]*",	OPTIONS+="last_rule"
KERNEL=="pty*",		OPTIONS+="last_rule"
SUBSYSTEM=="vc",	OPTIONS+="last_rule"

