# test/integration-tests/TEST-38-FREEZER/TEST-38-FREEZER.units/TEST-38-FREEZER-sleep.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Service]
ExecStart=sleep 3600
