# etc/udev/redhat/40-redhat.rules
# do not edit this file, it will be overwritten on update

KERNEL=="hd*[!0-9]", SUBSYSTEMS=="ide", ATTRS{media}=="floppy", SYMLINK+="floppy floppy-%k", OPTIONS+="ignore_remove"
KERNEL=="hd*[0-9]", ATTRS{media}=="floppy", ATTRS{removable}=="1", SYMLINK+="floppy-%k", OPTIONS+="ignore_remove"

KERNEL=="fw*", PROGRAM="fw_unit_symlinks.sh %k %n", SYMLINK+="$result"

ACTION=="add", KERNEL=="sg[0-9]*", BUS=="scsi", ATTRS{type}=="[36]", SYMLINK+="scanner scanner-%k", MODE="0660"

ACTION=="add", KERNEL=="sg[0-9]*", BUS=="scsi", ATTRS{type}=="8", SYMLINK+="changer changer-%k", MODE="0660", GROUP="disk"

ACTION=="add", SUBSYSTEM=="usb_device", PROGRAM="/bin/sh -c 'K=%k; K=$${K#usbdev}; printf bus/usb/%%03i/%%03i $${K%%%%.*} $${K#*.}'", NAME="%c", MODE="0644"

BUS=="usb", KERNEL=="sd*", ATTRS{bInterfaceClass}=="08", ATTRS{bInterfaceSubClass}=="04", GROUP="floppy", MODE="0660", 	SYMLINK+="floppy floppy-%k"

KERNEL=="fd[0-9]*",		SYMLINK+="floppy floppy-%k"

ACTION=="add", SUBSYSTEM=="scsi_device", ATTRS{type}=="0|7|14", \
        RUN+="/sbin/modprobe sd_mod"
ACTION=="add", SUBSYSTEM=="scsi_device", ATTRS{type}=="[45]", \
        RUN+="/sbin/modprobe sr_mod"

KERNEL=="vcs",                 OWNER="vcsa", GROUP="tty"
KERNEL=="vcs[0-9]*",           OWNER="vcsa", GROUP="tty"
KERNEL=="vcsa",                OWNER="vcsa", GROUP="tty"
KERNEL=="vcsa[0-9]*",          OWNER="vcsa", GROUP="tty"
KERNEL=="vcc/*",               OWNER="vcsa", GROUP="tty"
