# etc/udev/redhat/51-hotplug.rules
# do not call hotplug.d and dev.d for "drivers" and "module" events
SUBSYSTEM=="drivers",   GOTO="hotplug_end"
SUBSYSTEM=="module",    GOTO="hotplug_end"
ACTION=="add",           GOTO="hotplug_comp"
ACTION=="remove",        GOTO="hotplug_comp"

GOTO="hotplug_end"

LABEL="hotplug_comp"
# compatibility support for the obsolete hotplug.d and dev.d directories
ENV{UDEVD_EVENT}=="1",  RUN+="/lib/udev/udev_run_hotplugd"
RUN+="/lib/udev/udev_run_devd"


LABEL="hotplug_end"
