# test2/b.service
[Unit]
Description=B
Wants=f.service

[Service]
ExecStart=/bin/true
