# test2/f.service
[Unit]
Description=F

[Service]
ExecStart=/bin/true
