# test2/g.service
[Unit]
Description=G
Conflicts=e.service

[Service]
ExecStart=/bin/true
