# test/test-fstab-generator/test-19-swap-from-cmdline.expected/systemd-mkswap@dev-sdy2.service
# Automatically generated by systemd-fstab-generator

[Unit]
Description=Make Swap on %f
Documentation=man:systemd-mkswap@.service(8)
DefaultDependencies=no
BindsTo=%i.device
Conflicts=shutdown.target
After=%i.device
Before=shutdown.target dev-sdy2.swap

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/lib/systemd/systemd-makefs swap /dev/sdy2
TimeoutSec=0
