# test/TEST-63-PATH/TEST-63-PATH.units/test63.path
# SPDX-License-Identifier: LGPL-2.1-or-later
[Path]
PathExists=/tmp/test63
# Make the unit friendly to slower machines
TriggerLimitIntervalSec=10
TriggerLimitBurst=10
