# test/TEST-63-PATH/TEST-63-PATH.units/test63-issue-24577.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Requires=test63-issue-24577-dep.service
After=test63-issue-24577-dep.service

[Service]
Type=oneshot
ExecStart=bash -c 'sleep infinity'
