# extras/hid2hci/70-hid2hci.rules
# do not edit this file, it will be overwritten on update

ACTION!="add", GOTO="hid2hci_end"
SUBSYSTEM!="usb", GOTO="hid2hci_end"

# Variety of Dell Bluetooth devices - it looks like a bit of an odd rule,
# because it is matching on a mouse device that is self powered, but that
# is where a HID report needs to be sent to switch modes.
#
# Known supported devices: 413c:8154, 413c:8158, 413c:8162
ATTR{bInterfaceProtocol}=="02", ATTRS{idVendor}=="413c", ATTRS{bmAttributes}=="e0", \
    RUN+="hid2hci --method dell -v $attr{idVendor} -p $attr{idProduct} --mode hci"

ENV{DEVTYPE}!="usb_device", GOTO="hid2hci_end"

# Logitech devices
ATTR{idVendor}=="046d", ATTR{idProduct}=="c70[345abce]", RUN+="hid2hci --method logitech -v $attr{idVendor} -p $attr{idProduct} --mode hci"
ATTR{idVendor}=="046d", ATTR{idProduct}=="c71[34bc]", RUN+="hid2hci --method logitech -v $attr{idVendor} -p $attr{idProduct} --mode hci"

# CSR devices (in HID mode)
ATTR{idVendor}=="0a12", ATTR{idProduct}=="1000", RUN+="hid2hci --method csr -v $attr{idVendor} -p $attr{idProduct} --mode hci"
ATTR{idVendor}=="0458", ATTR{idProduct}=="1000", RUN+="hid2hci --method csr -v $attr{idVendor} -p $attr{idProduct} --mode hci"
ATTR{idVendor}=="05ac", ATTR{idProduct}=="1000", RUN+="hid2hci --method csr -v $attr{idVendor} -p $attr{idProduct} --mode hci"

LABEL="hid2hci_end"
