# test/test-fstab-generator/test-12-dev-sda.expected/systemd-fsck-root.service
# Automatically generated by systemd-fstab-generator

[Unit]
Description=File System Check on /dev/sda1
Documentation=man:systemd-fsck-root.service(8)
DefaultDependencies=no
BindsTo=dev-sda1.device
Conflicts=shutdown.target
After=initrd-root-device.target local-fs-pre.target dev-sda1.device
Before=shutdown.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/lib/systemd/systemd-fsck /dev/sda1
TimeoutSec=0
