# mkosi/mkosi.sanitizers/mkosi.extra/usr/lib/systemd/system/systemd-coredump@.service.d/10-asan.conf
# SPDX-License-Identifier: LGPL-2.1-or-later

[Service]
# systemd-coredump may call get_user_creds(), which may pull in instrumented
# systemd NSS modules and may trigger fatal LSAN error.
EnvironmentFile=-/usr/lib/systemd/systemd-asan-env
