# test/test-fstab-generator/test-21-clone.fstab.expected.sysroot/systemd-clone@myclone.service
# Automatically generated by systemd-fstab-generator

[Unit]
Description=Create dm-clone device myclone for /mnt/data
Documentation=man:dmsetup(8) man:fstab(5) man:systemd-fstab-generator(8)
DefaultDependencies=no
BindsTo=dev-sda1.device dev-sdb1.device dev-sdc1.device
Requires=dev-sda1.device dev-sdb1.device dev-sdc1.device
After=dev-sda1.device dev-sdb1.device dev-sdc1.device
Before=mnt-data.mount
Conflicts=shutdown.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/systemd-clone add 'myclone' '/dev/sda1' '/dev/sdb1' '/dev/sdc1' ''
ExecStop=/usr/bin/systemd-clone remove myclone
TimeoutSec=0
