# mkosi/mkosi.initrd.conf/mkosi.extra/usr/lib/systemd/system/user@.service.d/fdstore.conf
# For tests exercising the FD store we need the unit in the initrd to have these
# settings, or the fdstore content will be dropped in the initrd
[Service]
FileDescriptorStoreMax=20
FileDescriptorStorePreserve=yes
