# test/TEST-87-TIMER/TEST-87-TIMER.units/realtime-test.timer
[Unit]
Description=Testing systemd timers

[Timer]
OnCalendar=*:*:0/5
AccuracySec=1us
DeferReactivation=true

[Install]
WantedBy=timers.target
