# test/TEST-74-AUX-UTILS/TEST-74-AUX-UTILS.units/realtime-test.service
[Unit]
Description=Testing systemd timers

[Service]
Type=simple
ExecStart=sh -c 'date +%%s >>/tmp/realtime-test.log ; sleep 5'
