# test/units/b.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=B
Wants=f.service

[Service]
ExecStart=/bin/true
