# test/units/dml-discard.slice
[Unit]
Description=DML discard slice

[Slice]
DefaultMemoryLow=
