# test/units/g.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=G
Conflicts=e.service

[Service]
ExecStart=/bin/true
