# test/units/grandchild.service
[Unit]
Description=Grandchild Service

[Service]
Slice=parent-deep.slice
Type=oneshot
ExecStart=/bin/true
