# test/units/h.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=H
Wants=g.service

[Service]
ExecStart=/bin/true
