# test/units/i.service
[Unit]
Description=I
Conflicts=a.service d.service
Wants=b.service
After=b.service

[Service]
ExecStart=/bin/true
