# test/units/loopy3.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Service]
ExecStart=/bin/true

[Unit]
Conflicts=loopy4.service
