# test/units/parent-deep.slice
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Deeper Parent Slice

[Slice]
MemoryAccounting=yes
