# test/units/sched_rr_ok.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Default prio for RR

[Service]
ExecStart=true
CPUSchedulingPolicy=rr
