# test/units/sched_rr_bad.service
[Unit]
Description=Bad sched priority for RR

[Service]
ExecStart=/bin/true
CPUSchedulingPolicy=rr
CPUSchedulingPriority=0
CPUSchedulingPriority=100
