# test/units/sched_idle_bad.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Bad sched priority for Idle

[Service]
ExecStart=/bin/true
CPUSchedulingPriority=1
