# test/units/sched_idle_ok.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Sched idle with prio 0

[Service]
ExecStart=/bin/true
CPUSchedulingPriority=0
