# test/units/son.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Son Service

[Service]
Slice=parent.slice
Type=oneshot
ExecStart=true
CPUShares=100
