# test/units/testsuite-08.service
[Unit]
Description=TEST-08-ISSUE-2730

[Service]
ExecStartPre=rm -f /failed /testok
ExecStart=sh -x -c 'mount -o remount,rw /dev/vda1 && echo OK >/testok; systemctl poweroff'
Type=oneshot
