# test/units/testsuite-15.service
[Unit]
Description=TEST-15-DROPIN

[Service]
ExecStart=/usr/lib/systemd/tests/testdata/units/%N.sh
Type=oneshot
