# test/units/testsuite-80-fails-on-restart.target
[Unit]
Requires=testsuite-80-fails-on-restart.service
