# test/units/testsuite-77-server@.service
[Unit]
Description=TEST-77-OPENFILE server

[Service]
ExecStart=echo "Socket"
StandardInput=socket
StandardOutput=socket
