# test/units/testsuite-44.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=TESTSUITE-44-LOG-NAMESPACE
Before=getty-pre.target
Wants=getty-pre.target
Wants=systemd-journald@foobar.socket systemd-journald-varlink@foobar.socket
After=systemd-journald@foobar.socket systemd-journald-varlink@foobar.socket

[Service]
ExecStartPre=rm -f /failed /testok
ExecStart=/usr/lib/systemd/tests/testdata/units/%N.sh
Type=oneshot
