# test/units/testsuite-69.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=TEST-69-SHUTDOWN

[Service]
Type=oneshot
ExecStart=true
