# test/units/testsuite-51-repro-1.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Issue 16115 Repro with on-abnormal

[Service]
Type=simple
Restart=on-abnormal
ExecCondition=/bin/false
ExecStart=sleep 100
RestartSec=1
