# test/units/testsuite-57.service
[Unit]
Description=TEST-57-RELOADING-RESTART

[Service]
Type=oneshot
ExecStart=/usr/lib/systemd/tests/testdata/units/%N.sh
