# test/units/testsuite-25.service
[Unit]
Description=TEST-25-IMPORT

[Service]
ExecStart=/usr/lib/systemd/tests/testdata/units/%N.sh
Type=oneshot
