# test/units/testsuite-77-server.socket
[Unit]
Description=TEST-77-OPENFILE server socket

[Socket]
ListenStream=/tmp/test.sock
Accept=yes
