# test/units/testsuite-26.service
[Unit]
Description=TEST-26-SETENV

[Service]
ExecStart=/usr/lib/systemd/tests/testdata/units/%N.sh
Type=oneshot
