# test/units/TEST-66-DEVICE-ISOLATION-device-isolation.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Service that uses device isolation

[Service]
DevicePolicy=strict
DeviceAllow=/dev/null r
StandardOutput=file:/tmp/TEST-66-DEVICE-ISOLATION.serviceresults
ExecStartPre=rm -f /tmp/TEST-66-DEVICE-ISOLATION.serviceresults
ExecStart=bash -c "while true; do sleep 0.01 && echo meow >/dev/null && echo thisshouldnotbehere; done"
