# test/units/testsuite-40.service
[Unit]
Description=TEST-40-EXEC-COMMAND-EX

[Service]
ExecStart=/usr/lib/systemd/tests/testdata/units/%N.sh
Type=oneshot
