# test/units/testsuite-51.service
[Unit]
Description=TEST-51-ISSUE-16115

[Service]
ExecStartPre=rm -f /failed /testok
ExecStart=/usr/lib/systemd/tests/testdata/units/%N.sh
Type=oneshot
