# test/units/testsuite-58.service
[Unit]
Description=TEST-58-REPART

[Service]
ExecStart=/usr/lib/systemd/tests/testdata/units/%N.sh
Type=oneshot
