# test/units/testsuite-59.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=TEST-59-RELOADING-RESTART

[Service]
Type=oneshot
ExecStart=/usr/lib/systemd/tests/testdata/units/%N.sh
