# test/units/testsuite-56.service
[Unit]
Description=TEST-56-EXIT-TYPE

[Service]
ExecStart=/usr/lib/systemd/tests/testdata/units/%N.sh
Type=oneshot
