# test/units/testsuite-24.service
[Unit]
Description=TEST-24-UNIT-TESTS

[Service]
ExecStartPre=rm -f /failed /testok
ExecStart=/usr/lib/systemd/tests/testdata/units/%N.sh
Type=oneshot
