# test/units/testsuite-27.service
[Unit]
Description=TEST-27-STDOUTFILE

[Service]
ExecStart=/usr/lib/systemd/tests/testdata/units/%N.sh
Type=oneshot
