# test/units/testsuite-80-succeeds-on-restart.target
[Unit]
Requires=testsuite-80-succeeds-on-restart.service
