# test/units/testsuite-17.service
[Unit]
Description=TEST-17-UDEV-WANTS

[Service]
ExecStart=/usr/lib/systemd/tests/testdata/units/%N.sh
Type=oneshot
