# test/units/testsuite-57-success.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Succeeding unit
OnSuccess=testsuite-57-fail.service

[Service]
ExecStart=/bin/true
