# test/units/testsuite-99.service
[Unit]
Description=TEST-99-RESTART-MODE

[Service]
ExecStart=/usr/lib/systemd/tests/testdata/units/%N.sh
Type=oneshot
