# test/units/testsuite-59.service
[Unit]
Description=TEST-59-TYPE-ANY

[Service]
ExecStart=/usr/lib/systemd/tests/testdata/units/%N.sh
Type=oneshot
