# test/units/testsuite-80-succeeds-on-restart.service
[Service]
Type=oneshot
ExecStart=/usr/lib/systemd/tests/testdata/units/%N.sh
Restart=on-failure
