# test/units/testsuite-37.service
[Unit]
Description=TEST-37-RUNTIMEDIRECTORYPRESERVE

[Service]
ExecStart=/usr/lib/systemd/tests/testdata/units/%N.sh
Type=oneshot
