# test/units/testsuite-46.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=TEST-46-HOMED
Wants=getty-pre.target
Before=getty-pre.target
Wants=systemd-homed.service
After=systemd-homed.service

[Service]
ExecStartPre=rm -f /failed /testok
ExecStart=/usr/lib/systemd/tests/testdata/units/%N.sh
Type=oneshot
NotifyAccess=all
