# test/units/testsuite-34.service
[Unit]
Description=TEST-34-DYNAMICUSERMIGRATE

[Service]
ExecStart=/usr/lib/systemd/tests/testdata/units/%N.sh
Type=oneshot
