# test/units/testsuite-28.service
[Unit]
Description=TEST-28-PERCENTJ-WANTEDBY
# Testsuite: Ensure %j Wants directives work
Wants=specifier-j-wants.service
After=specifier-j-wants.service

[Service]
ExecStartPre=rm -f /failed /testok
ExecStart=true
Type=oneshot
