# test/units/testsuite-36.service
[Unit]
Description=TEST-36-NUMAPOLICY

[Service]
ExecStart=/usr/lib/systemd/tests/testdata/units/%N.sh
Type=oneshot
