# test/units/testsuite-62-3.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=TEST-62-RESTRICT-IFACES-deny-list
[Service]
ExecStart=/bin/sh -c '! ping -c 1 -W 0.2 192.168.113.1'
ExecStart=/bin/sh -c '! ping -c 1 -W 0.2 192.168.113.5'
ExecStart=/bin/sh -c 'ping -c 1 -W 0.2 192.168.113.9'
RestrictNetworkInterfaces=~veth0
RestrictNetworkInterfaces=~veth1
Type=oneshot
