# test/units/testsuite-47-repro.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Issue 14566 Repro

[Service]
ExecStart=/usr/lib/systemd/tests/testdata/units/%N.sh
ExecStopPost=/bin/true
KillMode=mixed
