# test/units/testsuite-23.service
[Unit]
Description=TEST-23-TYPE-EXEC

[Service]
ExecStart=/usr/lib/systemd/tests/testdata/units/%N.sh
Type=oneshot
