# test/units/testsuite-56-testchill.service
[Unit]
Description=No memory pressure

[Service]
Slice=testsuite-56-workload.slice
ExecStart=sleep infinity
