# test/units/testsuite-57-prop-stop-one.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Stop Propagation Receiver
Wants=testsuite-57-prop-stop-two.service
After=testsuite-57-prop-stop-two.service
StopPropagatedFrom=testsuite-57-prop-stop-two.service

[Service]
ExecStart=/bin/sleep infinity
ExecStopPost=systemctl kill --kill-whom=main -sUSR2 testsuite-57.service
