# test/units/testsuite-54.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=TESTSUITE-54-CREDS

[Service]
ExecStartPre=rm -f /failed /testok
ExecStart=/usr/lib/systemd/tests/testdata/units/%N.sh
Type=oneshot
