# test/units/testsuite-57-bound-by.service
[Unit]
Description=Unit with BoundBy=

[Service]
ExecStart=/bin/sleep 0.7
