# test/units/testsuite-60.service
[Unit]
Description=TEST-59-MOUNT-RATELIMIT

[Service]
Type=oneshot
ExecStart=/usr/lib/systemd/tests/testdata/units/%N.sh
