# test/units/testsuite-31.service
[Unit]
Description=TEST-31-DEVICE-ENUMERATION

[Service]
ExecStart=/usr/lib/systemd/tests/testdata/units/%N.sh
Type=oneshot
