# test/units/testsuite-59.service
[Unit]
Description=TEST-59-RELOADING-RESTART

[Service]
Type=oneshot
ExecStart=/usr/lib/systemd/tests/testdata/units/%N.sh
