# test/units/testsuite-47-repro.service
[Unit]
Description=Issue 14566 Repro

[Service]
ExecStart=/usr/lib/systemd/tests/testdata/units/%N.sh
ExecStopPost=/bin/true
KillMode=mixed
