# test/units/testsuite-39.service
[Unit]
Description=TEST-39-EXECRELOAD

[Service]
ExecStart=/usr/lib/systemd/tests/testdata/units/%N.sh
Type=oneshot
