# test/units/testsuite-24.service
[Unit]
Description=TEST-24-CRYPTSETUP
After=multi-user.target

[Service]
ExecStartPre=rm -f /failed /testok
ExecStart=sh -x -e -c 'mountpoint /var; systemctl --state=failed --no-legend --no-pager >/failed; echo OK >/testok'
Type=oneshot
