# test/units/testsuite-30.service
[Unit]
Description=TEST-30-ONCLOCKCHANGE

[Service]
ExecStartPre=rm -f /failed /testok
ExecStart=/usr/lib/systemd/tests/testdata/units/%N.sh
Type=oneshot
