# test/units/testsuite-55-testmunch.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Create some memory pressure

[Service]
MemoryHigh=12M
Slice=testsuite-55-workload.slice
ExecStart=stress --timeout 3m --vm 10 --vm-bytes 200M --vm-keep --vm-stride 1
