# test/units/testsuite-47.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=TEST-47-ISSUE-14566

[Service]
ExecStartPre=rm -f /failed /testok
ExecStart=/usr/lib/systemd/tests/testdata/units/%N.sh
Type=oneshot
