# test/units/testsuite-21.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Fuzz our D-Bus interfaces with dfuzzer
After=dbus.service multi-user.target
Wants=dbus.service multi-user.target

[Service]
ExecStartPre=rm -f /failed /skipped /testok
ExecStart=/usr/lib/systemd/tests/testdata/units/%N.sh
Type=oneshot
