# test/units/testsuite-57-prop-stop-two.service
[Unit]
Description=Stop Propagation Sender

[Service]
ExecStart=/bin/sleep 1.5
