# test/units/testsuite-55-testbloat.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Create a lot of memory pressure

[Service]
# A VERY small memory.high will cause the 'stress' (trying to use a lot of memory)
# to throttle and be put under heavy pressure.
MemoryHigh=3M
Slice=testsuite-55-workload.slice
ExecStart=stress --timeout 3m --vm 10 --vm-bytes 200M --vm-keep --vm-stride 1
