# test/units/testsuite-57-fail.service
[Unit]
Description=Failing unit
OnFailure=testsuite-57-uphold.service

[Service]
ExecStart=/bin/false
