# test/units/testsuite-19.service
[Unit]
Description=TEST-19-DELEGATE

[Service]
ExecStart=/usr/lib/systemd/tests/testdata/units/%N.sh
Type=oneshot
