# test/units/unit-with-multiple-dashes.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=A unit with multiple dashes
Documentation=man:test

[Service]
ExecStart=/bin/true
