# etc/udev/suse/80-sysconfig.rules
ENV{SYSCONFIG}=="no|off|0", GOTO="sysconfig_end"

# SUSE system integration
SUBSYSTEM=="pci", ACTION=="add", RUN+="/sbin/hwup bus-pci-%k -o hotplug"
SUBSYSTEM=="usb", ACTION=="add", RUN+="/sbin/hwup usb-devpath-%p -o hotplug"
SUBSYSTEM=="ieee1394", ACTION=="add", RUN+="/sbin/hwup ieee1394-devpath-%p -o hotplug"
SUBSYSTEM=="scsi", DEVPATH=="/devices/css*", ACTION=="remove", RUN+="/sbin/hwdown scsi-devpath-%p -o hotplug"
SUBSYSTEM=="scsi_host", ACTION=="add", SYSFS{proc_name}=="zfcp", RUN+="/sbin/hwup scsi_host-devpath-%p -o hotplug"
SUBSYSTEM=="vio", ACTION=="add", RUN+="/sbin/hwup bus-macio-%k -o hotplug"
SUBSYSTEM=="macio", ACTION=="add", RUN+="/sbin/hwup bus-macio-%k -o hotplug"
SUBSYSTEM=="ccw", ACTION=="add", RUN+="/sbin/hwup ccw-devpath-%p -o hotplug"
SUBSYSTEM=="ccwgroup", ACTION=="add", RUN+="/sbin/hwup ccwgroup-devpath-%p -o hotplug"
SUBSYSTEM=="pnp", ACTION=="add", RUN+="/sbin/hwup pnp-devpath-%p -o hotplug"
SUBSYSTEM=="ide", ACTION=="add", RUN+="/sbin/hwup ide-devpath-%p -o hotplug"
SUBSYSTEM=="input", KERNEL=="input[0-9]*", ACTION=="add", RUN+="/sbin/hwup input-devpath-%p -o hotplug"

LABEL="sysconfig_end"
