# etc/udev/suse/90-hal.rules
RUN+="socket:/org/freedesktop/hal/udev_event"
