# etc/udev/suse/90-hal.rules
# send all events to HAL

RUN+="socket:/org/freedesktop/hal/udev_event"
