# etc/udev/suse/90-hal.rules
# send all events to HAL
RUN+="/usr/sbin/hal.hotplug"

