# etc/udev/suse/95-udev-late.rules
# event to be catched by udevmonitor
RUN+="socket:/org/kernel/udev/monitor"

