# etc/udev/suse/95-udev-late.rules
# debug events to /events, but not on default udev_log="err"
ENV{UDEV_LOG}=="[4-9]", RUN+="eventrecorder.sh"

# event to be catched by udevmonitor
RUN+="socket:/org/kernel/udev/monitor"

