# test/integration-tests/TEST-03-JOBS/TEST-03-JOBS.units/restart-on-failure.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Service]
Type=oneshot
ExecStart=touch /tmp/test-03-restart-failure-flag
