# test/integration-tests/TEST-03-JOBS/TEST-03-JOBS.units/unstoppable.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=echo "I'm unstoppable!"
ExecStop=systemctl start --no-block unstoppable.service
