# test/TEST-04-SECCOMP/will-not-fail2.service
[Unit]
Description=Reset SystemCallFilter

[Service]
ExecStart=/bin/echo "Foo bar"
SystemCallFilter=
