# test/testsuite-23.units/testsuite-23-upheldby-install.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Unit that sets UpheldBy= through [Install]

[Service]
ExecStart=sleep infinity

[Install]
UpheldBy=testsuite-23-retry-uphold.service
