# test/testsuite-23.units/testsuite-23-bound-by.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Unit with BoundBy=

[Service]
ExecStart=/bin/sleep 0.7
