# test/testsuite-23.units/testsuite-23-joins-namespace-of-6.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
JoinsNamespaceOf=testsuite-23-joins-namespace-of-8.service

[Service]
Type=notify
NotifyAccess=all
MountAPIVFS=yes
PrivateTmp=yes
ExecStart=bash -c 'touch /tmp/shared-private-file-x && systemd-notify --ready && sleep infinity'
