# test/testsuite-23.units/testsuite-23.MaxInstances@.service
[Unit]
MaxInstances=2

[Service]
ExecStart=true
RemainAfterExit=yes
