# test/testsuite-23.units/testsuite-23-prop-stop-one.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Stop Propagation Receiver
Wants=testsuite-23-prop-stop-two.service
After=testsuite-23-prop-stop-two.service
StopPropagatedFrom=testsuite-23-prop-stop-two.service

[Service]
ExecStart=sleep infinity
ExecStopPost=systemctl kill --kill-whom=main -sUSR2 testsuite-23.service
