# test/testsuite-23.units/testsuite-23-success.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Succeeding unit
OnSuccess=testsuite-23-fail.service

[Service]
ExecStart=true
