# test/testsuite-23.units/testsuite-23.MaxUnits@.service
[Service]
ExecStart=true
RemainAfterExit=yes
Slice=testsuite-23.MaxUnits.slice
