# test/test-cgroup-mask/grandchild.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Grandchild Service

[Service]
Slice=parent-deep.slice
Type=oneshot
ExecStart=true
