# etc/udev/rules.d/95-udev-late.rules
# do not edit this file, it will be overwritten on update

# event to be catched by udevmonitor
RUN+="socket:/org/kernel/udev/monitor"

