# test/test-fstab-generator/test-20-swap-from-cmdline.expected/systemd-mkswap@dev-sdy3.service
# Automatically generated by systemd-fstab-generator

[Unit]
Description=Make Swap on %f
Documentation=man:systemd-mkswap@.service(8)

DefaultDependencies=no
BindsTo=%i.device
After=%i.device
Before=dev-sdy3.swap
Conflicts=shutdown.target
Before=shutdown.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/lib/systemd/systemd-makefs swap /dev/sdy3
TimeoutSec=infinity
