# test/TEST-23-UNIT-FILE.units/TEST-23-UNIT-FILE-bound-by.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Unit with BoundBy=

[Service]
ExecStart=sleep 0.7
