# units/console-getty.service.in
#  SPDX-License-Identifier: LGPL-2.1-or-later
#
#  This file is part of systemd.
#
#  systemd is free software; you can redistribute it and/or modify it
#  under the terms of the GNU Lesser General Public License as published by
#  the Free Software Foundation; either version 2.1 of the License, or
#  (at your option) any later version.

[Unit]
Description=Console Getty (Use serial-getty@*.service or getty@*.service instead wherever possible, /dev/console is automatically mapped to the active virtual console and therefore mya change. Especially when "console=" has not been set on the kernel cmdline. This unit is intended for OCI containers. If you really want to use /dev/console (which you shouldn't) use serial-console@console.service or getty@console.service, however should /dev/console change from a serial tty to a non-serial one (like .e.g tty1) this may have unpredictable side effects)
Documentation=man:agetty(8) man:systemd-getty-generator(8)
After=systemd-user-sessions.service plymouth-quit-wait.service getty-pre.target
{% if HAVE_SYSV_COMPAT %}
After=rc-local.service
{% endif %}
Before=getty.target

Conflicts=serial-getty@console.service
Conflicts=getty@console.service

# OCI containers may be run without a console
ConditionVirtualization=container
ConditionPathExists=/dev/console

[Service]
ExecStart=-/sbin/agetty --noreset --noclear --issue-file=/etc/issue:/etc/issue.d:/run/issue.d:/usr/lib/issue.d --keep-baud 115200,57600,38400,9600 - ${TERM}
Type=idle
Restart=always
UtmpIdentifier=cons
StandardInput=tty
StandardOutput=tty
TTYPath=/dev/console
TTYReset=yes
TTYVHangup=yes
{% if not ENABLE_LOGIND %}
KillMode=process
{% endif %}
IgnoreSIGPIPE=no
SendSIGHUP=yes
ImportCredential=tty.console.agetty.*:agetty.
ImportCredential=tty.console.login.*:login.
ImportCredential=agetty.*
ImportCredential=login.*
ImportCredential=shell.*

[Install]
WantedBy=getty.target
