# units/systemd-vconsole-setup.service.in
#  SPDX-License-Identifier: LGPL-2.1-or-later
#
#  This file is part of systemd.
#
#  systemd is free software; you can redistribute it and/or modify it
#  under the terms of the GNU Lesser General Public License as published by
#  the Free Software Foundation; either version 2.1 of the License, or
#  (at your option) any later version.

{% if ENABLE_FIRSTBOOT %}
# Note that this service exists in two copies: systemd-vconsole-setup.service
# and systemd-vconsole-setup-reload.service. The second one is used internally
# by systemd-firstboot.service to start a second copy of the unit without
# ordering dependencies. The first copy should be used in all other cases.
{% endif %}

[Unit]
Description=Virtual Console Setup
Documentation=man:systemd-vconsole-setup.service(8) man:vconsole.conf(5)

ConditionPathExists=/dev/tty0

DefaultDependencies=no
Before=sysinit.target
Before=initrd-switch-root.target shutdown.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart={{ROOTLIBEXECDIR}}/systemd-vconsole-setup
LoadCredential=vconsole.keymap
LoadCredential=vconsole.keymap_toggle
LoadCredential=vconsole.font
LoadCredential=vconsole.font_map
LoadCredential=vconsole.font_unimap
