# units/systemd-logger.service.in
[Unit]
Description=systemd Logging Daemon

[Service]
ExecStart=@pkglibexecdir@/systemd-logger
Type=simple
