# units/systemd-tmpfiles-clean.service.in
#  This file is part of systemd.
#
#  systemd is free software; you can redistribute it and/or modify it
#  under the terms of the GNU General Public License as published by
#  the Free Software Foundation; either version 2 of the License, or
#  (at your option) any later version.

[Unit]
Description=Cleanup of Temporary Directories
DefaultDependencies=no
Wants=local-fs.target
After=systemd-readahead-collect.service systemd-readahead-replay.service local-fs.target
Before=basic.target shutdown.target
ConditionPathExists=/etc/tmpfiles.d

[Service]
Type=oneshot
ExecStart=@rootlibexecdir@/systemd-tmpfiles --clean
IOSchedulingClass=idle
