# units/systemd-udev-trigger.service.in
[Unit]
Description=udev Coldplug all Devices
Wants=systemd-udev.service
After=systemd-udev-kernel.socket systemd-udev-control.socket
DefaultDependencies=no
ConditionCapability=CAP_MKNOD

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@bindir@/udevadm trigger --type=subsystems --action=add ; @bindir@/udevadm trigger --type=devices --action=add
