# units/systemd-readahead-done.service.in
#  This file is part of systemd.
#
#  systemd is free software; you can redistribute it and/or modify it
#  under the terms of the GNU Lesser General Public License as published by
#  the Free Software Foundation; either version 2.1 of the License, or
#  (at your option) any later version.

[Unit]
Description=Stop Read-Ahead Data Collection
Documentation=man:systemd-readahead-replay.service(8)
DefaultDependencies=no
Conflicts=shutdown.target
After=default.target
Before=shutdown.target
ConditionVirtualization=no

[Service]
Type=oneshot
ExecStart=@SYSTEMD_NOTIFY@ --readahead=done

[Install]
Also=systemd-readahead-collect.service
