# units/systemd-initctl.socket
[Unit]
Description=systemd /dev/initctl Compatibility

[Socket]
ListenFIFO=/dev/initctl
SocketMode=0600
