# units/systemd-udev-settle.service.in
#  This file is part of systemd.
#
#  systemd is free software; you can redistribute it and/or modify it
#  under the terms of the GNU Lesser General Public License as published by
#  the Free Software Foundation; either version 2.1 of the License, or
#  (at your option) any later version.

# This service is usually not enabled by default. If enabled, it
# acts as a barrier for basic.target -- so all later services will
# wait for udev completely finishing its coldplug run.
#
# If needed, to work around broken or non-hotplug-aware services,
# it might be enabled unconditionally, or pulled-in on-demand by
# the services that assume a fully populated /dev at startup. It
# should not be used or pulled-in ever on systems without such
# legacy services running.

[Unit]
Description=udev Wait for Complete Device Initialization
Documentation=man:udev(7) man:systemd-udevd.service(8)
DefaultDependencies=no
Wants=systemd-udevd.service
After=systemd-udev-trigger.service
Before=basic.target
ConditionCapability=CAP_MKNOD

[Service]
Type=oneshot
TimeoutSec=180
RemainAfterExit=yes
ExecStart=@bindir@/udevadm settle
