# units/systemd-journalctl.socket
[Unit]
Description=Journal Log Access Socket
Documentation=man:journalctl(1)
DefaultDependencies=no
Before=sockets.target

[Socket]
ListenStream=/run/systemd/io.systemd.Journalctl
Symlinks=/run/varlink/registry/io.systemd.Journalctl
FileDescriptorName=varlink
SocketMode=0666
Accept=yes
MaxConnectionsPerSource=16
