# units/systemd-network-generator.service.in
#  SPDX-License-Identifier: LGPL-2.1+
#
#  This file is part of systemd.
#
#  systemd is free software; you can redistribute it and/or modify it
#  under the terms of the GNU Lesser General Public License as published by
#  the Free Software Foundation; either version 2.1 of the License, or
#  (at your option) any later version.

[Unit]
Description=Generate network units from Kernel command line
DefaultDependencies=no
Before=network-pre.target

[Service]
Type=oneshot
ExecStart=@rootlibexecdir@/systemd-network-generator

[Install]
WantedBy=network-pre.target
