# units/systemd-imds-metrics@.service.in
#  SPDX-License-Identifier: LGPL-2.1-or-later
#
#  This file is part of systemd.
#
#  systemd is free software; you can redistribute it and/or modify it
#  under the terms of the GNU Lesser General Public License as published by
#  the Free Software Foundation; either version 2.1 of the License, or
#  (at your option) any later version.

[Unit]
Description=Cloud Instance Metadata Metrics
Documentation=man:systemd-imds(1)
DefaultDependencies=no
Conflicts=shutdown.target initrd-switch-root.target
Before=shutdown.target initrd-switch-root.target
After=systemd-imdsd.socket

[Service]
# Runs systemd-imds as an io.systemd.Metrics provider. It connects to
# systemd-imdsd to acquire the actual data. Runs as root so the daemon's Polkit
# check on io.systemd.InstanceMetadata.Get() passes.
ExecStart={{LIBEXECDIR}}/systemd-imds
LockPersonality=yes
MemoryDenyWriteExecute=yes
NoNewPrivileges=yes
