# units/systemd-tpm2-swtpm.service.in
#  SPDX-License-Identifier: LGPL-2.1-or-later
#
#  This file is part of systemd.
#
#  systemd is free software; you can redistribute it and/or modify it
#  under the terms of the GNU Lesser General Public License as published by
#  the Free Software Foundation; either version 2.1 of the License, or
#  (at your option) any later version.

[Unit]
Description=Fallback Software TPM
Documentation=man:systemd-tpm2-swtpm.service(8)
DefaultDependencies=no
After=systemd-sysusers.service
Wants=modprobe@tpm_vtpm_proxy.service
After=modprobe@tpm_vtpm_proxy.service
After=sysefi.mount
Wants=sysefi.mount
Before=tpm2.target sysinit.target

[Service]
Type=notify
RemainAfterExit=yes
ExecStart={{LIBEXECDIR}}/systemd-tpm2-swtpm
ExecStop=swtpm_ioctl --unix %t/systemd/swtpm/socket -v
ExecStop=swtpm_ioctl --unix %t/systemd/swtpm/socket -s
RuntimeDirectory=systemd/swtpm

[Install]
WantedBy=factory-reset.target
