# units/systemd-random-seed-load.service.in
#  This file is part of systemd.
#
#  systemd is free software; you can redistribute it and/or modify it
#  under the terms of the GNU Lesser General Public License as published by
#  the Free Software Foundation; either version 2.1 of the License, or
#  (at your option) any later version.

[Unit]
Description=Load Random Seed
DefaultDependencies=no
Conflicts=shutdown.target
After=systemd-readahead-collect.service systemd-readahead-replay.service
Before=sysinit.target shutdown.target
RequiresMountsFor=@RANDOM_SEED@

[Service]
Type=oneshot
ExecStart=@rootlibexecdir@/systemd-random-seed load
