# units/systemd-journal-upload.service.in
#  This file is part of systemd.
#
#  systemd is free software; you can redistribute it and/or modify it
#  under the terms of the GNU Lesser General Public License as published by
#  the Free Software Foundation; either version 2.1 of the License, or
#  (at your option) any later version.

[Unit]
Description=Journal Remote Upload Service
After=network.target

[Service]
ExecStart=@rootlibexecdir@/systemd-journal-upload \
          --save-state
User=systemd-journal-upload
PrivateTmp=yes
PrivateDevices=yes
WatchdogSec=20min

# If there are many split up journal files we need a lot of fds to
# access them all and combine
LimitNOFILE=16384

[Install]
WantedBy=multi-user.target
