# units/systemd-time-wait-sync.service.in
#  SPDX-License-Identifier: LGPL-2.1-or-later
#
#  This file is part of systemd.
#
#  systemd is free software; you can redistribute it and/or modify it
#  under the terms of the GNU Lesser General Public License as published by
#  the Free Software Foundation; either version 2.1 of the License, or
#  (at your option) any later version.

[Unit]
Description=Wait Until Kernel Time Synchronized
Documentation=man:systemd-time-wait-sync.service(8)

# Note that this tool doesn't need CAP_SYS_TIME itself, but its primary
# usecase is to run in conjunction with a local NTP service such as
# systemd-timesyncd.service, which is conditioned this way. There might be
# niche usecases where running this service independently is desired, but let's
# make this all "just work" for the general case, and leave it to local
# modifications to make it work in the remaining cases.

ConditionCapability=CAP_SYS_TIME
ConditionVirtualization=!container

DefaultDependencies=no
Before=time-sync.target shutdown.target
Wants=time-sync.target
Conflicts=shutdown.target

[Service]
Type=oneshot
ExecStart={{ROOTLIBEXECDIR}}/systemd-time-wait-sync
TimeoutStartSec=infinity
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
