# units/systemd-initctl.service.in
[Meta]
Description=systemd /dev/initctl Compatibility

[Service]
ExecStart=@libexecdir@/systemd-initctl
Type=simple
