# units/systemd-readahead-replay.service.in
#  This file is part of systemd.
#
#  systemd is free software; you can redistribute it and/or modify it
#  under the terms of the GNU Lesser General Public License as published by
#  the Free Software Foundation; either version 2.1 of the License, or
#  (at your option) any later version.

[Unit]
Description=Replay Read-Ahead Data
Documentation=man:systemd-readahead-replay.service(8)
DefaultDependencies=no
Conflicts=shutdown.target
Before=sysinit.target shutdown.target
ConditionPathExists=!/run/systemd/readahead/noreplay
ConditionPathExists=/.readahead
ConditionVirtualization=no

[Service]
Type=notify
ExecStart=@rootlibexecdir@/systemd-readahead replay
RemainAfterExit=yes
StandardOutput=null
OOMScoreAdjust=1000

[Install]
WantedBy=default.target
