# units/systemd-logger.socket
[Unit]
Description=systemd Logging Socket

[Socket]
ListenStream=@/org/freedesktop/systemd1/logger
