# units/systemd-growfs-root.service.in
#  SPDX-License-Identifier: LGPL-2.1-or-later
#
#  This file is part of systemd.
#
#  systemd is free software; you can redistribute it and/or modify it
#  under the terms of the GNU Lesser General Public License as published by
#  the Free Software Foundation; either version 2.1 of the License, or
#  (at your option) any later version.

[Unit]
Description=Grow Root File System
Documentation=man:systemd-growfs-root.service(8)

DefaultDependencies=no
# Partitions must be resized first
After=systemd-repart.service
# Make sure the root fs is marked writable before we resize it
After=systemd-remount-fs.service
# The file system should be writable before root-fs-ready.target
Before=root-fs-ready.target

Conflicts=shutdown.target
Before=shutdown.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart={{ROOTLIBEXECDIR}}/systemd-growfs /
TimeoutSec=0
