# units/systemd-update-utmp-shutdown.service.in
#  This file is part of systemd.
#
#  systemd is free software; you can redistribute it and/or modify it
#  under the terms of the GNU Lesser General Public License as published by
#  the Free Software Foundation; either version 2.1 of the License, or
#  (at your option) any later version.

[Unit]
Description=Update UTMP about System Shutdown
DefaultDependencies=no
After=local-fs.target sysinit.target auditd.service systemd-update-utmp-runlevel.service
Before=poweroff.service reboot.service halt.service

[Service]
Type=oneshot
ExecStart=@rootlibexecdir@/systemd-update-utmp shutdown
