# units/systemd-initctl.service.in
[Unit]
Description=systemd /dev/initctl Compatibility

[Service]
ExecStart=@pkglibexecdir@/systemd-initctl
Type=simple
