# units/systemd-bootchart.service.in
#  This file is part of systemd.
#
#  systemd is free software; you can redistribute it and/or modify it
#  under the terms of the GNU Lesser General Public License as published by
#  the Free Software Foundation; either version 2.1 of the License, or
#  (at your option) any later version.

# Note: it's usually a better idea to run systemd-bootchart via the
# init= kernel command line switch. See the man page for details.

[Unit]
Description=Boot Process Profiler
Documentation=man:systemd-bootchart.service(1) man:bootchart.conf(5)
DefaultDependencies=no

[Service]
ExecStart=@rootlibexecdir@/systemd-bootchart -r

[Install]
WantedBy=sysinit.target
