# units/systemd-udev.service.in
[Unit]
Description=udev Kernel Device Manager
Wants=systemd-udev-control.socket systemd-udev-kernel.socket
After=systemd-udev-control.socket systemd-udev-kernel.socket
Before=basic.target
DefaultDependencies=no
ConditionCapability=CAP_MKNOD

[Service]
Type=notify
OOMScoreAdjust=-1000
Sockets=systemd-udev-control.socket systemd-udev-kernel.socket
Restart=on-failure
ExecStart=@rootlibexecdir@/systemd-udevd
