# units/systemd-journalctl@.service
[Unit]
Description=Journal Log Access Service
Documentation=man:journalctl(1)
DefaultDependencies=no
Conflicts=shutdown.target
Before=shutdown.target

[Service]
ExecStart=journalctl
