# units/systemd-logger.service.in
[Meta]
Description=systemd Logging Daemon

[Service]
ExecStart=@libexecdir@/systemd-logger
Type=simple
