# rules.d/70-amd-pmf.rules.in
# do not edit this file, it will be overwritten on update

ACTION!="change", GOTO="pmf_end"
DRIVERS!="amd-pmf", GOTO="pmf_end"

ENV{EVENT_ID}=="2", RUN+="{{BINDIR}}/loginctl lock-sessions"

LABEL="pmf_end"
