# test/test-path/path-mycustomunit.service
[Unit]
Description=Service Test Path Unit

[Service]
ExecStart=sleep infinity
Type=exec
RemainAfterExit=true
