# test/test-path/path-existsglob.path
# SPDX-License-Identifier: LGPL-2.1-or-later

[Unit]
Description=Test PathExistsGlob

[Path]
PathExistsGlob=/tmp/test-path_existsglob*

[Install]
WantedBy=multi-user.target
