# test/test-path/path-modified.path
[Unit]
Description=Test PathModified

[Path]
PathModified=/tmp/test-path_modified

[Install]
WantedBy=multi-user.target
