# test/test-path/path-exists.path
# SPDX-License-Identifier: LGPL-2.1-or-later

[Unit]
Description=Test PathExists

[Path]
PathExists=/tmp/test-path_exists

[Install]
WantedBy=multi-user.target
