# test/test-path/path-service.service
[Unit]
Description=Service Test for Path units

[Service]
ExecStart=/bin/true
Type=oneshot
