# test/test-path/path-makedirectory.service
[Unit]
Description=Service Test for Path units

[Service]
ExecStart=sleep infinity
Type=exec
RemainAfterExit=true
