# test/test-path/path-deactivationtoggle-directorynotempty.path
# SPDX-License-Identifier: LGPL-2.1-or-later

[Unit]
Description=Test DeactivationToggle with DirectoryNotEmpty

[Path]
DirectoryNotEmpty=/tmp/test-path_deactivationtoggle/
Unit=path-deactivationtoggle.service
DeactivationToggle=true

[Install]
WantedBy=multi-user.target
