# test/testsuite-52.units/testsuite-52.service
[Unit]
Description=Testsuite service

[Service]
ExecStart=/usr/lib/systemd/tests/testdata/%N.units/%N.sh
Type=oneshot
