# test/integration-tests/TEST-23-UNIT-FILE/TEST-23-UNIT-FILE.units/success-failure-test-success.service
[Service]
Type=notify
ExecStart=bash -c "echo success >> /tmp/success-failure-test-result && systemd-notify --ready && sleep infinity"
