# test/integration-tests/TEST-23-UNIT-FILE/TEST-23-UNIT-FILE.units/TEST-23-UNIT-FILE-openfile-server@.service
[Unit]
Description=OpenFile= test server service

[Service]
ExecStart=echo "Socket"
StandardInput=socket
StandardOutput=socket
