# test/integration-tests/TEST-23-UNIT-FILE/TEST-23-UNIT-FILE.units/TEST-23-UNIT-FILE-openfile-server.socket
[Unit]
Description=OpenFile= test server socket

[Socket]
ListenStream=/tmp/test.sock
Accept=yes
