# test/integration-tests/TEST-23-UNIT-FILE/TEST-23-UNIT-FILE.units/TEST-23-UNIT-FILE-non-namespaced.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Service]
RuntimeMaxSec=5
Type=notify
RemainAfterExit=yes
ExecStart=sh -c 'systemd-notify --ready; until grep -q -F MARKER_RUNTIME /tmp/testfile-marker-runtime; do sleep 0.1; done; exit 0'
