# mkosi.images/system/mkosi.extra/usr/lib/systemd/system/iscsid.socket.d/99-skip-in-container.conf
# SPDX-License-Identifier: LGPL-2.1-or-later
# CI runners may be containers running iscsi without a separate IPC namespace
# so starting this unit can fail with
# (@ISCSIADM_ABSTRACT_NAMESPACE): Address already in use
# so this unit is disabled in containers.
[Unit]
ConditionVirtualization=!container
