# test/test-fstab-generator/test-19-mounts-from-cmdline.expected/hoge-withx20space.mount
# Automatically generated by systemd-fstab-generator

[Unit]
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
SourcePath=/proc/cmdline
Before=remote-fs.target

[Mount]
What=//foo￾bar
Where=/hoge/with space
Type=cifs
Options=rw,seclabel
