# test/test-fstab-generator/test-19-mounts-from-cmdline.expected/hoge-without_options.mount
# Automatically generated by systemd-fstab-generator

[Unit]
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
SourcePath=/proc/cmdline
Before=local-fs.target
Requires=systemd-fsck@dev-sdx2.service
After=systemd-fsck@dev-sdx2.service
After=blockdev@dev-sdx2.target

[Mount]
What=/dev/sdx2
Where=/hoge/without_options
