# test/test-fstab-generator/test-19-mounts-from-cmdline.expected/sysroot-foo-also_in_initrd.mount
# Automatically generated by systemd-fstab-generator

[Unit]
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
SourcePath=/proc/cmdline
Before=initrd-fs.target
Requires=systemd-fsck@dev-sdx6.service
After=systemd-fsck@dev-sdx6.service
After=blockdev@dev-sdx6.target

[Mount]
What=/dev/sdx6
# Canonicalized from /foo/also_in_initrd
Where=/sysroot/foo/also_in_initrd
Options=x-initrd.mount
