# test/test-fstab-generator/test-19-mounts-from-cmdline.expected/systemd-fsck-usr.service
# Automatically generated by systemd-fstab-generator

[Unit]
Description=File System Check on /dev/sdx5
Documentation=man:systemd-fsck-usr.service(8)

DefaultDependencies=no
BindsTo=dev-sdx5.device
Conflicts=shutdown.target
After=local-fs-pre.target dev-sdx5.device
Before=shutdown.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/lib/systemd/systemd-fsck /dev/sdx5
TimeoutSec=infinity
ImportCredential=fsck.*
