# test/testsuite-28.units/testsuite-28-pre.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Service]
ExecStart=rm -f /failed /testok
Type=oneshot
