# test/testsuite-28.units/testsuite-28-pre.service
[Service]
ExecStart=rm -f /failed /testok
Type=oneshot
