# rules/90-hda-powersave.rules
# do not edit this file, it will be overwritten on update

ACTION=="remove", GOTO="hda_powersave_end"

SUBSYSTEM=="sound", DRIVERS=="snd_hda_intel", \
  IMPORT{builtin}="hwdb 'snd_hda_intel:$attr{subsystem_vendor}:$attr{subsystem_device}:$attr{[dmi/id]modalias}'"

ENV{SND_HDA_INTEL_POWERSAVE}!="", \
  RUN+="/bin/sh -c 'echo $env{SND_HDA_INTEL_POWERSAVE} > /sys/module/snd_hda_intel/parameters/power_save'"

LABEL="hda_powersave_end"
