# test/integration-tests/TEST-07-PID1/TEST-07-PID1.units/defer-reactivation.service
[Unit]
Description=Test for DeferReactivation=yes (service)

[Service]
Type=simple
ExecStart=sh -c 'date +%%s >>/tmp/defer-reactivation.log; sleep 5'
