# test/integration-tests/TEST-07-PID1/TEST-07-PID1.units/prefix-shell.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Service]
Type=oneshot
Environment=SHLVL=100
ExecStartPre=-|false
ExecStart=|@logger -p crit "with login shell $$SHELL: lvl $$SHLVL"
ExecStart=:|"str='with normal shell'" logger -p crit "$$str"
ExecStart=|echo YAY! >/tmp/TEST-07-PID1.prefix-shell.flag
