# test/test-fstab-generator/test-16-tmpfs.expected/sysroot.mount
# Automatically generated by systemd-fstab-generator

[Unit]
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
SourcePath=/proc/cmdline
Before=initrd-root-fs.target

[Mount]
What=rootfs
Where=/sysroot
Type=tmpfs
Options=rw,mode=0755
