# units/user@.service.d/10-defaults.conf
[Unit]
After=systemd-user-sessions.service user-runtime-dir@%i.service dbus.service
