# mkosi.images/app-reload/mkosi.extra/usr/lib/systemd/system/foo.service
[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=echo foo

[Install]
WantedBy=multi-user.target
