# test/testsuite-03.units/fails-on-restart.target
[Unit]
Requires=fails-on-restart.service
