# test/testsuite-03.units/fails-on-restart-restartdirect.target
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Requires=fails-on-restart-restartdirect.service
