# test/testsuite-03.units/fails-on-restart-restartquick.target
[Unit]
Requires=fails-on-restart-restartquick.service
