# test/testsuite-03.units/propagatestopto-only.target
[Unit]
PropagatesStopTo=sleep-infinity-simple.service
