# test/testsuite-03.units/propagatestopto-and-pullin.target
[Unit]
Requires=sleep-infinity-simple.service
After=sleep-infinity-simple.service
PropagatesStopTo=sleep-infinity-simple.service
