# test/testsuite-03.units/sleep.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Sleep for 1 minute

[Service]
Type=oneshot
ExecStart=/bin/sleep 60
