# test/testsuite-03.units/sleep-infinity-simple.service
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Sleep infinitely

[Service]
Type=simple
ExecStart=/bin/sleep infinity
