# test/testsuite-03.units/succeeds-on-restart-restartquick.target
[Unit]
Requires=succeeds-on-restart-restartquick.service
