# test/testsuite-03.units/succeeds-on-restart-restartquick.service
[Service]
Type=oneshot
ExecStart=/usr/lib/systemd/tests/testdata/testsuite-03.units/succeeds-on-restart.sh
Restart=on-failure
RestartQuick=yes
