# test/testsuite-03.units/succeeds-on-restart.target
[Unit]
Requires=succeeds-on-restart.service
