# test/testsuite-03.units/succeeds-on-restart.service
[Service]
Type=oneshot
ExecStart=/usr/lib/systemd/tests/testdata/testsuite-03.units/%N.sh
Restart=on-failure
RestartQuick=yes
