# test/testsuite-08.units/-.mount
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Before=local-fs.target

[Mount]
What=/dev/sda1
Where=/
Options=noatime

[Install]
WantedBy=local-fs.target
Alias=root.mount
