# extras/mtd_probe/75-probe_mtd.rules.in
# do not edit this file, it will be overwritten on update

ACTION!="add", GOTO="mtd_probe_end"

KERNEL=="mtd*ro", IMPORT{program}="mtd_probe $tempnode"
KERNEL=="mtd*ro", ENV{MTD_FTL}=="smartmedia", RUN+="@sbindir@/modprobe -bv sm_ftl"

LABEL="mtd_probe_end"
