# extras/mtd_probe/75-probe_mtd.rules
#probe mtd devices for FTL

ACTION!="add", GOTO="mtd_probe_end"

KERNEL=="mtd*ro", IMPORT{program}="mtd_probe $tempnode"
KERNEL=="mtd*ro", ENV{MTD_FTL}=="smartmedia", RUN+="/sbin/modprobe sm_ftl"


LABEL="mtd_probe_end"

